# comms/smstools/pkg/smsd.rc
#!/bin/ksh
#
# $OpenBSD: smsd.rc,v 1.3 2018/01/11 19:27:01 rpe Exp $

daemon="${TRUEPREFIX}/bin/smsd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o _smsd /var/run/smsd
}

rc_cmd $1
