# comms/smstools/pkg/smsd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/smsd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o _smsd /var/run/smsd
}

rc_cmd $1
