# net/miniupnp/miniupnpd/pkg/miniupnpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/miniupnpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
