# sysutils/consul/pkg/consul.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/consul agent"
daemon_flags="-config-dir ${SYSCONFDIR}/consul.d"
daemon_user="_consul"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_stop_signal=INT

rc_cmd $1
