# net/dhcpcd/pkg/dhcpcd.rc
#!/bin/ksh
#
# $OpenBSD: dhcpcd.rc,v 1.3 2018/03/29 19:31:22 sthen Exp $

daemon="${TRUEPREFIX}/sbin/dhcpcd"
daemon_flags="-Mq"

. /etc/rc.d/rc.subr

pexp="dhcpcd: .*"

rc_cmd $1
