# net/dhcpcd/pkg/dhcpcd.rc
#!/bin/ksh
#
# $OpenBSD: dhcpcd.rc,v 1.2 2018/01/11 19:27:05 rpe Exp $

daemon="${TRUEPREFIX}/sbin/dhcpcd"
daemon_flags="-Mq"

. /etc/rc.d/rc.subr

rc_cmd $1
