# net/dhcpcd/pkg/dhcpcd.rc
#!/bin/ksh
#
# $OpenBSD: dhcpcd.rc,v 1.4 2020/04/02 16:33:19 sthen Exp $

daemon="${TRUEPREFIX}/sbin/dhcpcd"
daemon_flags="-Mq"

. /etc/rc.d/rc.subr

pexp="dhcpcd:.*master.*"

rc_cmd $1
