# net/ntopng/pkg/ntopng.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/ntopng --daemon"

. /etc/rc.d/rc.subr

rc_pre() {
	/usr/bin/install -d -o _ntopng -g wheel -m 750 /var/run/ntopng
}

rc_cmd $1
