# x11/dbus/pkg/messagebus.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/dbus-daemon --system"

. /etc/rc.d/rc.subr

rc_pre() {
	rm -f /var/run/dbus/pid
	install -d -o _dbus -g _dbus /var/run/dbus
}

rc_cmd $1
