# mail/dovecot/pkg/dovecot.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/dovecot"

. /etc/rc.d/rc.subr

rc_cmd $1
