# mail/courier-imap/pkg/courier_pop3.rc
#!/bin/sh
#
# $OpenBSD: courier_pop3.rc,v 1.3 2010/12/27 14:50:23 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/pop3d.rc"

. /etc/rc.d/rc.subr

pexp=".*courierlogger.*pop3d.pid"

rc_start() {
	${daemon} start 2>/dev/null
}

rc_stop() {
	${daemon} stop 2>/dev/null
}

rc_cmd $1
