# mail/courier-imap/pkg/courier_imap.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/imapd.rc"

. /etc/rc.d/rc.subr

pexp="${LOCALBASE}/sbin/courierlogger .* -start -name=imapd .*"

rc_start() {
	${rcexec} "${daemon} start"
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
