# mail/courier-imap/pkg/courier_imap.rc
#!/bin/sh
#
# $OpenBSD: courier_imap.rc,v 1.6 2015/07/15 21:14:06 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/imapd.rc"

. /etc/rc.d/rc.subr

pexp="${LOCALBASE}/sbin/courierlogger .* -start -name=imapd .*"

rc_start() {
	${rcexec} "${daemon} start"
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
