# mail/milter-spamd/pkg/milter_spamd.rc
#!/bin/sh
#
# $OpenBSD: milter_spamd.rc,v 1.2 2010/12/11 12:57:37 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/libexec/milter-spamd"
rc_reload=NO

rc_cmd $1
