# mail/milter-spamd/pkg/milter_spamd.rc
#!/bin/ksh
#
# $OpenBSD: milter_spamd.rc,v 1.7 2018/01/11 19:27:03 rpe Exp $

daemon="${TRUEPREFIX}/libexec/milter-spamd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
