# mail/milter-spamd/pkg/milter_spamd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/milter-spamd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
