# mail/milter-spamd/pkg/milter_spamd.rc
#!/bin/sh
#
# $OpenBSD: milter_spamd.rc,v 1.5 2012/12/05 09:38:41 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/milter-spamd"

. /etc/rc.d/rc.subr

pexp="${daemon}"
rc_reload=NO

rc_cmd $1
