# net/gmid/pkg/gmid.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/gmid"

. /etc/rc.d/rc.subr

rc_configtest() {
	${daemon} -n ${daemon_flags}
}

rc_cmd $1
