# www/varnish/pkg/varnishd.rc
#!/bin/sh
#
# $OpenBSD: varnishd.rc,v 1.3 2016/08/23 13:48:23 gonzalo Exp $

daemon="${PREFIX}/sbin/varnishd"
daemon_flags="-j unix,user=_varnish,ccgroup=_varnish -f ${SYSCONFDIR}/varnish/default.vcl"

. /etc/rc.d/rc.subr

pexp="varnishd: Varnish-Mgr $(hostname)"
rc_reload=NO

rc_cmd $1
