# net/olsrd/pkg/PLIST-main
@conflict olsrd-plugins-*
@pkgpath net/olsrd
@pkgpath net/olsrd,-plugins
bin/olsr-topology-view.pl
lib/olsrd-plugins/
lib/olsrd-plugins/olsrd_dot_draw.so
lib/olsrd-plugins/olsrd_httpinfo.so
lib/olsrd-plugins/olsrd_jsoninfo.so
lib/olsrd-plugins/olsrd_mini.so
lib/olsrd-plugins/olsrd_nameservice.so
lib/olsrd-plugins/olsrd_netjson.so
lib/olsrd-plugins/olsrd_pgraph.so
lib/olsrd-plugins/olsrd_secure.so
lib/olsrd-plugins/olsrd_txtinfo.so
lib/olsrd-plugins/olsrd_watchdog.so
@man man/man5/olsrd.conf.5
@man man/man8/olsrd.8
@bin sbin/olsrd
share/examples/olsrd/
@sample ${SYSCONFDIR}/olsrd/
share/examples/olsrd/olsrd.conf
@sample ${SYSCONFDIR}/olsrd/olsrd.conf
share/examples/olsrd/olsrd.conf.default.lq
share/examples/olsrd/olsrd.conf.default.rfc
share/examples/olsrd/olsrd.conf.default.txt
@rcscript /etc/rc.d/olsrd
