# net/lldpd/pkg/lldpd.rc
#!/bin/ksh
#
# $OpenBSD: lldpd.rc,v 1.5 2018/01/11 19:27:06 rpe Exp $

daemon="${TRUEPREFIX}/sbin/lldpd"

. /etc/rc.d/rc.subr

pexp="lldpd: monitor."
rc_reload=NO

rc_pre() {
	install -d -o _lldpd /var/run/lldpd
}

rc_cmd $1
