# audio/spotifyd/pkg/spotifyd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/spotifyd"
daemon_user="_spotifyd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
