# misc/gpsd/pkg/gpsd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/gpsd"

. /etc/rc.d/rc.subr

rc_cmd $1
