# misc/gpsd/pkg/gpsd.rc
#!/bin/ksh
#
# $OpenBSD: gpsd.rc,v 1.1 2021/08/07 12:41:15 sthen Exp $

daemon="${TRUEPREFIX}/sbin/gpsd"

. /etc/rc.d/rc.subr

rc_cmd $1
