# www/kibana/pkg/kibana.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/kibana/bin/kibana"
daemon_user="_kibana"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/bin/node ${TRUEPREFIX}/kibana/bin/../src/cli/dist${daemon_flags:+ ${daemon_flags}}"
rc_bg=YES
rc_reload=NO

rc_cmd $1
