# www/kibana/pkg/kibana.rc
#!/bin/sh
#
# $OpenBSD: kibana.rc,v 1.3 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/kibana/bin/kibana"
daemon_user="_kibana"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/bin/node ${TRUEPREFIX}/kibana/bin/../src/cli${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
