# net/thingsd/pkg/thingsd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/thingsd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
