# sysutils/alertmanager/pkg/alertmanager.rc
#!/bin/ksh
#
# $OpenBSD: alertmanager.rc,v 1.1.1.1 2022/02/16 18:26:58 denis Exp $

daemon="${TRUEPREFIX}/bin/alertmanager"
daemon_flags="--config.file ${SYSCONFDIR}/alertmanager/alertmanager.yml"
daemon_flags="${daemon_flags} --storage.path '${LOCALSTATEDIR}/alertmanager'"
daemon_logger="daemon.info"
daemon_user="_alertmanager"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1

