# net/dnscrypt-proxy/pkg/dnscrypt_proxy.rc
#!/bin/sh
#
# $OpenBSD: dnscrypt_proxy.rc,v 1.3 2016/02/04 14:29:25 sthen Exp $

daemon="${TRUEPREFIX}/sbin/dnscrypt-proxy -d --user=_dnscrypt-proxy"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
