# net/dnscrypt-proxy/pkg/dnscrypt_proxy.rc
#!/bin/ksh
#
# $OpenBSD: dnscrypt_proxy.rc,v 1.4 2018/01/11 19:27:05 rpe Exp $

daemon="${TRUEPREFIX}/sbin/dnscrypt-proxy -d --user=_dnscrypt-proxy"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
