# net/dnscrypt-proxy/pkg/dnscrypt_proxy.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/dnscrypt-proxy"
daemon_flags="-config ${SYSCONFDIR}/dnscrypt-proxy.toml"

. /etc/rc.d/rc.subr

pexp="${daemon}${daemon_flags:+ ${daemon_flags}}.*"

rc_bg=YES
rc_reload=NO

rc_cmd $1
