# databases/puppetdb/pkg/puppetdb.rc
#!/bin/sh
#
# $OpenBSD: puppetdb.rc,v 1.3 2014/07/09 14:24:38 jasper Exp $

daemon=$(${LOCALBASE}/bin/javaPathHelper -c puppetdb)
daemon_flags="-cp ${PREFIX}/share/puppetdb/puppetdb.jar clojure.main -m com.puppetlabs.puppetdb.core services -c ${SYSCONFDIR}/puppetdb/conf.d/puppetdb.ini"
daemon_timeout=120

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
