# sysutils/collectd/pkg/collectd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/collectd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	chgrp _collectd /etc/collectd.conf
	chmod g+r /etc/collectd.conf
}

rc_cmd $1
