# geo/osrm-backend/pkg/osrm_routed.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/osrm-routed"
daemon_user="_osrm_routed"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
