# comms/zigbee2mqtt/pkg/zigbee2mqtt.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/zigbee2mqtt"
daemon_user="_z2m"
daemon_logger="daemon.info"

. /etc/rc.d/rc.subr

pexp="node ${daemon}"
rc_reload=NO
rc_bg=YES

rc_cmd $1
