# sysutils/tenshi/pkg/tenshi.rc
#!/bin/ksh
#
# $OpenBSD: tenshi.rc,v 1.2 2021/11/09 07:39:08 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/tenshi"

. /etc/rc.d/rc.subr

pexp="perl: tenshi"

rc_cmd $1
