# x11/slim/pkg/slim.rc
#!/bin/sh
#
# $OpenBSD: slim.rc,v 1.2 2011/03/09 18:51:49 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/slim"
daemon_flags="-d"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	(sleep 5; ${rcexec} "${daemon} ${daemon_flags}")&
}

rc_cmd $1
