# x11/slim/pkg/slim.rc
#!/bin/sh
#
# $OpenBSD: slim.rc,v 1.11 2017/07/21 23:06:08 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/slim -d"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
