# x11/slim/pkg/slim.rc
#!/bin/ksh
#
# $OpenBSD: slim.rc,v 1.12 2018/01/11 19:27:12 rpe Exp $

daemon="${TRUEPREFIX}/bin/slim -d"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
