# www/puppetboard/pkg/puppetboard.rc
#!/bin/sh
#
# $OpenBSD: puppetboard.rc,v 1.2 2014/06/17 11:55:11 jasper Exp $

daemon="puppetboard_wsgi"
daemon_flags="--workers 4"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN}.*${daemon}:application"

rc_start() {
	${rcexec} "cd ${INSTDIR} && \
		${LOCALBASE}/bin/gunicorn --daemon ${daemon_flags} ${daemon}:application"
}

rc_cmd $1
