# mail/sympa/pkg/sympa.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/sympa"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	rm -f /var/spool/lock/sympa
	${rcexec} "${daemon} start"
}

rc_stop() {
	${daemon} stop
}

rc_check() {
	${daemon} status
}

rc_cmd $1
