# mail/sympa/pkg/sympa.rc
#!/bin/ksh
#
# $OpenBSD: sympa.rc,v 1.7 2018/01/11 19:27:04 rpe Exp $

daemon="${TRUEPREFIX}/sbin/sympa"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	rm -f /var/spool/lock/sympa
	${rcexec} "${daemon} start"
}

rc_stop() {
	${daemon} stop
}

rc_check() {
	${daemon} status
}

rc_cmd $1
