# mail/sympa/pkg/sympa.rc
#!/bin/sh
#
# $OpenBSD: sympa.rc,v 1.3 2011/03/09 18:51:48 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/sympa"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	${daemon} start
}

rc_stop() {
	${daemon} stop
}

rc_check() {
	${daemon} status
}

rc_cmd $1
