# audio/hgd/pkg/hgd_netd.rc
#!/bin/sh
#
# $OpenBSD: hgd_netd.rc,v 1.2 2011/07/25 23:07:49 edd Exp $

daemon="${TRUEPREFIX}/sbin/hgd-netd"
daemon_user="_hgd"

. /etc/rc.d/rc.subr

rc_cmd $1
