# net/dnsdist/pkg/dnsdist.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/dnsdist --supervised"
daemon_flags="-u _dnsdist"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
