# net/ngircd/pkg/ngircd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/ngircd"

. /etc/rc.d/rc.subr

rc_cmd $1
