# security/suricata/pkg/suricata.rc
#!/bin/ksh
#
# $OpenBSD: suricata.rc,v 1.2 2018/02/15 16:19:41 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/suricata -D"

. /etc/rc.d/rc.subr

rc_cmd $1
