# mail/rspamd/pkg/rspamd.rc
#!/bin/sh
#
# $OpenBSD: rspamd.rc,v 1.2 2016/01/16 10:39:59 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/rspamd"
daemon_flags="-u _rspamd -g _rspamd"

. /etc/rc.d/rc.subr

pexp="rspamd: main process"

rc_pre() {
	install -d /var/run/rspamd
}

rc_cmd $1
