# sysutils/facette/pkg/facette.rc
#!/bin/ksh
#
# $OpenBSD: facette.rc,v 1.5 2018/08/31 14:04:40 danj Exp $

daemon="${TRUEPREFIX}/bin/facette"
daemon_user="_facette"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
