# sysutils/facette/pkg/facette.rc
#!/bin/sh
#
# $OpenBSD: facette.rc,v 1.3 2017/05/27 19:57:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/facette"
daemon_user="_facette"

# uncomment to enable debugging
#daemon_flags="-L debug -l ${LOCALSTATEDIR}/facette/log"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
