# sysutils/facette/pkg/facette.rc
#!/bin/sh
#
# $OpenBSD: facette.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/facette"
daemon_user="_facette"

# uncomment to enable debugging
#daemon_flags="-L debug -l ${LOCALSTATEDIR}/facette/log"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
