# net/ipcad/pkg/ipcad.rc
#!/bin/sh
#
# $OpenBSD: ipcad.rc,v 1.1 2011/06/06 08:34:57 sthen Exp $

daemon="${TRUEPREFIX}/sbin/ipcad"
daemon_flags="-d"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
