# mail/grommunio/gromox/pkg/gromox_http.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/gromox/http"
daemon_logger="mail.info"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_pre() {
	install -d -o _gromox -g _gromox /var/db/gromox
}

rc_cmd $1
