# mail/grommunio/gromox/pkg/gromox_midb.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/gromox/midb"
daemon_logger="mail.info"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_cmd $1
