# mail/grommunio/gromox/pkg/gromox_event.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/gromox/event"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_cmd $1
