# mail/grommunio/gromox/pkg/gromox_delivery.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/gromox/delivery"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_pre() {
	install -d -o _gromox -g _gromox /var/run/gromox
	install -d -o _gromox -g _gromox /var/db/gromox/queue/mess
	install -d -o _gromox -g _gromox /var/db/gromox/queue/save
	install -d -o _gromox -g _gromox /var/db/gromox/queue/cache
}

rc_cmd $1
