# net/prosody/pkg/prosody.rc
#!/bin/ksh
#
# $OpenBSD: prosody.rc,v 1.5 2020/06/02 07:22:46 landry Exp $

daemon="${TRUEPREFIX}/sbin/prosody -D"
daemon_user="_prosody"

. /etc/rc.d/rc.subr

pexp="${MODLUA_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_usercheck=NO

rc_check() {
	# return 1 is needed: prosodyctl returns a higher error code on failure
	${TRUEPREFIX}/sbin/prosodyctl status || return 1
}

rc_cmd $1
