# net/fastnetmon/pkg/fastnetmon.rc
#!/bin/ksh
#
# $OpenBSD: fastnetmon.rc,v 1.2 2018/07/11 15:24:56 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/fastnetmon --daemonize"
daemon_user="_fastnetmon"

. /etc/rc.d/rc.subr

rc_cmd $1
