# telephony/resiprocate/pkg/reTurnServer.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/reTurnServer"
daemon_flags="${SYSCONFDIR}/reTurn/reTurnServer.config"

. /etc/rc.d/rc.subr

rc_cmd $1
