# sysutils/google-compute-engine/pkg/google_accounts_daemon.rc
#!/bin/ksh
#
# $OpenBSD: google_accounts_daemon.rc,v 1.2 2018/01/11 19:27:10 rpe Exp $

daemon="${MODPY_BIN} ${TRUEPREFIX}/bin/google_accounts_daemon"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
