# www/gitea/pkg/gitea.rc
#!/bin/ksh
#
# $OpenBSD: gitea.rc,v 1.6 2021/01/17 21:07:41 pvk Exp $

daemon="${TRUEPREFIX}/sbin/gitea"
daemon_user="_gitea"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
