# sysutils/consul-template/pkg/consul_template.rc
#!/bin/ksh
#
# $OpenBSD: consul_template.rc,v 1.5 2021/11/06 10:38:40 ajacoutot Exp $

daemon="${PREFIX}/bin/consul-template"
daemon_flags="-config ${SYSCONFDIR}/consul-template.d/default.conf"
daemon_user="_consul-template"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_stop_signal=INT

rc_cmd $1
