# devel/py-buildbot/buildbot/pkg/buildbot.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/buildbot"
daemon_flags="${HOMEDIR}"
daemon_user="_buildbot"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} .*logfile=twistd.log --python=buildbot.tac"

rc_start() {
	${rcexec} "${daemon} start ${daemon_flags}"
}

rc_reload() {
	${daemon} reconfig ${daemon_flags}
}

rc_stop() {
	${daemon} stop ${daemon_flags}
}

rc_cmd $1
