# net/sniproxy/pkg/sniproxy.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/sniproxy"
daemon_flags="-c ${SYSCONFDIR}/sniproxy.conf"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
