# net/minbif/pkg/minbif.rc
#!/bin/ksh
#
# $OpenBSD: minbif.rc,v 1.2 2018/01/11 19:27:06 rpe Exp $

daemon="${TRUEPREFIX}/bin/minbif -m 2"
daemon_flags="${SYSCONFDIR}/minbif.conf"
daemon_user="_minbif"

. /etc/rc.d/rc.subr

rc_cmd $1
