# telephony/stuntman/pkg/stunserver.rc
#!/bin/sh
#
# $OpenBSD: stunserver.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/stunserver"
daemon_user="_stuntman"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
