# mail/postgrey/pkg/postgrey.rc
#!/bin/ksh
#
# $OpenBSD: postgrey.rc,v 1.4 2018/01/11 19:27:04 rpe Exp $

daemon="${TRUEPREFIX}/libexec/postgrey -d"
daemon_flags="--inet=10023"

. /etc/rc.d/rc.subr

pexp="perl: ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
