# net/isc-dhcp/pkg/isc_dhcpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/dhcpd"
daemon_flags="-user _isc-dhcp -group _isc-dhcp"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	touch /var/db/dhcpd.leases
}

rc_cmd $1
