# sysutils/toad/pkg/toadd.rc
#!/bin/ksh
#
# $OpenBSD: toadd.rc,v 1.4 2021/11/09 07:39:08 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/toadd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
