# sysutils/toad/pkg/toadd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/toadd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
