# telephony/siproxd/pkg/siproxd.rc
#!/bin/ksh
#
# $OpenBSD: siproxd.rc,v 1.3 2018/01/11 19:27:11 rpe Exp $

daemon="${TRUEPREFIX}/sbin/siproxd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
