# www/apache-httpd/pkg/httpd2.rc
#!/bin/sh
#
# $OpenBSD: httpd2.rc,v 1.1 2012/07/18 08:45:42 sthen Exp $

daemon="${TRUEPREFIX}/sbin/httpd2"

. /etc/rc.d/rc.subr

rc_pre() {
	rm -f /var/apache2/logs/{ssl_mutex,httpd.lock,accept.lock}.*
}

rc_cmd $1
