# databases/openldap23/pkg/slapd.rc
#!/bin/sh
#
# $OpenBSD: slapd.rc,v 1.2 2011/01/07 12:26:50 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/slapd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	/usr/bin/install -d -o _openldap /var/run/openldap
}

rc_cmd $1
