# net/bitcoin/pkg/bitcoind.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/bitcoind"
daemon_timeout="300"
daemon_user="_bitcoind"
data_dir="${VARBASE}/bitcoin"
daemon_flags="-daemon -conf=/etc/bitcoin.conf -datadir=${data_dir}"

. /etc/rc.d/rc.subr

rc_reload=NO


rc_cmd $1
