# net/dictd/pkg/dictd.rc
#!/bin/ksh

daemon="/usr/local/sbin/dictd"
daemon_args="--listen-to localhost"

. /etc/rc.d/rc.subr

pexp="${daemon##*/} .*"

rc_cmd $1
