# security/heimdal/pkg/ipropd_slave.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/heimdal/libexec/ipropd-slave --detach"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
