# security/heimdal/pkg/ipropd_master.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/heimdal/libexec/ipropd-master --detach"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
