# security/heimdal/pkg/kadmind.rc
#!/bin/sh
#
# $OpenBSD: kadmind.rc,v 1.1 2014/07/13 14:10:14 ajacoutot Exp $

daemon="${TRUEPREFIX}/heimdal/libexec/kadmind"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
