# security/heimdal/pkg/kdc.rc
#!/bin/ksh
#
# $OpenBSD: kdc.rc,v 1.2 2018/01/11 19:27:09 rpe Exp $

daemon="${TRUEPREFIX}/heimdal/libexec/kdc --detach"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
