# security/heimdal/pkg/kdc.rc
#!/bin/sh
#
# $OpenBSD: kdc.rc,v 1.1 2014/07/13 14:10:14 ajacoutot Exp $

daemon="${TRUEPREFIX}/heimdal/libexec/kdc --detach"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
