# www/shiori/pkg/shiori.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/shiori server"
daemon_user="_shiori"
daemon_flags="-a localhost"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
