# net/miniupnp/minissdpd/pkg/minissdpd.rc
#!/bin/sh
#
# $OpenBSD: minissdpd.rc,v 1.2 2013/08/09 14:31:54 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/minissdpd"
daemon_flags="-i 127.0.0.1"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
