# net/miniupnp/minissdpd/pkg/minissdpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/minissdpd"
daemon_flags="-i 127.0.0.1"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
