# www/squid/stable/pkg/squid.rc
#!/bin/sh
#
# $OpenBSD: squid.rc,v 1.2 2014/04/08 12:22:48 sthen Exp $

daemon="${TRUEPREFIX}/sbin/squid"

. /etc/rc.d/rc.subr

rc_stop() {
	${daemon} ${daemon_flags} -k shutdown
}

rc_reload() {
	${daemon} ${daemon_flags} -k reconfigure
}

rc_cmd $1
