# mail/smtp-vilter/pkg/smtp_vilter.rc
#!/bin/sh
#
# $OpenBSD: smtp_vilter.rc,v 1.5 2010/12/27 14:50:23 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/smtp-vilter"

. /etc/rc.d/rc.subr

pexp="smtp-vilter:"

rc_cmd $1
