# mail/smtp-vilter/pkg/smtp_vilter.rc
#!/bin/ksh
#
# $OpenBSD: smtp_vilter.rc,v 1.7 2018/01/11 19:27:04 rpe Exp $

daemon="${TRUEPREFIX}/sbin/smtp-vilter"

. /etc/rc.d/rc.subr

pexp="smtp-vilter: parent"

rc_cmd $1
