# mail/smtp-vilter/pkg/smtp_vilter.rc
#!/bin/sh
#
# $OpenBSD: smtp_vilter.rc,v 1.4 2010/12/24 10:40:06 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/smtp-vilter"

rc_conf

pexp="smtp-vilter:"

rc_cmd $1
