# mail/smtp-vilter/pkg/smtp_vilter.rc
#!/bin/sh
#
# $OpenBSD: smtp_vilter.rc,v 1.3 2010/11/04 12:01:42 sthen Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/smtp-vilter"
pexp="smtp-vilter:"

rc_cmd $1
