# mail/smtp-vilter/pkg/smtp_vilter.rc
#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/smtp-vilter"
pexp="smtp-vilter:

rc_cmd $1
