# mail/smtp-vilter/pkg/smtp_vilter.rc
#!/bin/sh
#
# $OpenBSD: smtp_vilter.rc,v 1.2 2010/10/29 12:54:44 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/smtp-vilter"
pexp="smtp-vilter:

rc_cmd $1
