# mail/smtp-vilter/pkg/smtp_vilter.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/smtp-vilter"

. /etc/rc.d/rc.subr

pexp="smtp-vilter: parent"

rc_cmd $1
