# telephony/baresip/restund/pkg/restund.rc
#!/bin/ksh
#
# $OpenBSD: restund.rc,v 1.2 2018/01/11 19:27:11 rpe Exp $

daemon="${TRUEPREFIX}/sbin/restund"
daemon_user="_restund"

. /etc/rc.d/rc.subr

rc_cmd $1
