# telephony/baresip/restund/pkg/restund.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/restund"
daemon_user="_restund"

. /etc/rc.d/rc.subr

rc_cmd $1
