# net/rsync/pkg/rsyncd.rc
#!/bin/sh
#
# $OpenBSD: rsyncd.rc,v 1.2 2016/04/06 21:17:29 sthen Exp $

daemon="${TRUEPREFIX}/bin/rsync --daemon"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
