# net/rsync/pkg/rsyncd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/rsync --daemon"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
