# net/nagios/nsca-ng/pkg/nsca_ng.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/nsca-ng"
daemon_flags="-c ${SYSCONFDIR}/nsca-ng.cfg"

. /etc/rc.d/rc.subr

rc_cmd $1
