# sysutils/burp/pkg/burp.rc
#!/bin/ksh
#
# $OpenBSD: burp.rc,v 1.3 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/sbin/burp"
daemon_flags="-c ${SYSCONFDIR}/burp/burp-server.conf"

. /etc/rc.d/rc.subr

rc_pre() {
	chown _burp ${SYSCONFDIR}/burp/ssl_cert_ca.pem
	install -d -o _burp -g _burp /var/run/burp
}

rc_cmd $1
