# sysutils/nut/pkg/upsd.rc
#!/bin/ksh
#
# $OpenBSD: upsd.rc,v 1.3 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/sbin/upsd"

. /etc/rc.d/rc.subr

upsdrvctl="${TRUEPREFIX}/sbin/upsdrvctl"

rc_pre() {
	${upsdrvctl} start
}

rc_post() {
	${upsdrvctl} stop
}

rc_cmd $1
