# sysutils/nut/pkg/upsmon.rc
#!/bin/ksh
#
# $OpenBSD: upsmon.rc,v 1.2 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/sbin/upsmon"

. /etc/rc.d/rc.subr

rc_reload() {
	${daemon} -c reload
}

rc_stop() {
	${daemon} -c stop
}

rc_cmd $1
