# net/ptpd/pkg/ptpd.rc
#!/bin/sh
#
# $OpenBSD: ptpd.rc,v 1.1 2017/09/09 21:22:49 naddy Exp $

daemon="${TRUEPREFIX}/sbin/ptpd"

. /etc/rc.d/rc.subr

rc_cmd $1
