# net/ptpd/pkg/ptpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/ptpd"

. /etc/rc.d/rc.subr

rc_cmd $1
