# net/i2pd/pkg/i2pd.rc
#!/bin/ksh
#
# $OpenBSD: i2pd.rc,v 1.3 2021/11/09 07:39:08 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/i2pd --daemon"
daemon_user="_i2pd"
daemon_flags="--service --datadir=${LOCALSTATEDIR}/lib/i2pd --conf=${SYSCONFDIR}/i2pd/i2pd.conf --tunconf=${SYSCONFDIR}/i2pd/tunnels.conf --tunnelsdir=${SYSCONFDIR}/i2pd/tunnels.d"

. /etc/rc.d/rc.subr

rc_cmd $1
