# net/mosquitto/pkg/mosquitto.rc
#!/bin/ksh
#
# $OpenBSD: mosquitto.rc,v 1.2 2018/01/11 19:27:06 rpe Exp $

daemon="${TRUEPREFIX}/sbin/mosquitto -d"
daemon_flags="-c ${SYSCONFDIR}/mosquitto/mosquitto.conf"

. /etc/rc.d/rc.subr

rc_cmd $1
