# net/ushare/pkg/ushare.rc
#!/bin/sh
#
# $OpenBSD: ushare.rc,v 1.2 2010/11/03 10:55:30 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/ushare"
daemon_flags="-D"

rc_cmd $1
