# net/ushare/pkg/ushare.rc
#!/bin/sh
#
# $OpenBSD: ushare.rc,v 1.4 2010/12/27 14:50:23 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/ushare"
daemon_flags="-D"

. /etc/rc.d/rc.subr

rc_cmd $1
