# sysutils/reaction/pkg/reaction.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/reaction"
daemon_flags="start -c ${SYSCONFDIR}/reaction.conf"
daemon_execdir="/var/reaction"
daemon_logger="daemon.info"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_bg=YES

rc_cmd $1
