# net/sslh/pkg/sslh_fork.rc
#!/bin/sh
#
# $OpenBSD: sslh_fork.rc,v 1.2 2012/04/13 14:21:45 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/sslh-fork"
daemon_flags="--listen 0.0.0.0:443 --ssh 127.0.0.1:22 --ssl 127.0.0.1:443"

. /etc/rc.d/rc.subr

pexp="${daemon}"
rc_reload=NO

rc_pre() {
	install -d -o _sslh /var/run/sslh
}

rc_cmd $1
