# databases/puppetdb5/pkg/puppetdb.rc
#!/bin/ksh

daemon=$(${LOCALBASE}/bin/javaPathHelper -c puppetdb)
daemon_flags="-cp ${PREFIX}/share/puppetdb/puppetdb.jar clojure.main -m puppetlabs.puppetdb.main --config ${SYSCONFDIR}/puppetlabs/puppetdb/conf.d/ -b ${SYSCONFDIR}/puppetlabs/puppetdb/bootstrap.cfg"
daemon_user="_puppetdb"
daemon_timeout=120

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
