# net/netbsd-iscsi-target/pkg/iscsi_target.rc
#!/bin/sh

daemon="${TRUEPREFIX}/bin/iscsi-target"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
