# security/cyrus-sasl2/pkg/saslauthd.rc
#!/bin/sh
#
# $OpenBSD: saslauthd.rc,v 1.2 2010/10/29 12:54:44 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/saslauthd"
daemon_flags="-a getpwent"

rc_cmd $1
