# security/cyrus-sasl2/pkg/saslauthd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/saslauthd"
daemon_flags="-a getpwent"

. /etc/rc.d/rc.subr

rc_cmd $1
