# security/cyrus-sasl2/pkg/saslauthd.rc
#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/saslauthd"
daemon_flags="-a getpwent"

rc_cmd $1
