# www/c-icap/c-icap/pkg/c_icap.rc
#!/bin/ksh
#
# $OpenBSD: c_icap.rc,v 1.2 2018/01/11 19:27:11 rpe Exp $

daemon="${TRUEPREFIX}/bin/c-icap"

. /etc/rc.d/rc.subr

rc_pre() {
	install -d -o _c-icap /var/run/c-icap
}

rc_cmd $1
