# graphics/sane-backends/pkg/saned.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/saned -a"
daemon_flags="_saned"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
