# graphics/sane-backends/pkg/saned.rc
#!/bin/sh
#
# $OpenBSD: saned.rc,v 1.7 2010/12/27 14:50:22 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/saned"
daemon_flags="-a _saned"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
