# graphics/sane-backends/pkg/saned.rc
#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/libexec/saned"
daemon_flags="-a _saned"

rc_reload() {
	rc_err "$0: reload is not supported"
}

rc_cmd $1
