# www/nginx/devel/pkg/enginx.rc
#!/bin/sh
#
# $OpenBSD: enginx.rc,v 1.1 2012/09/11 16:01:45 sthen Exp $

daemon="${TRUEPREFIX}/sbin/nginx"

. /etc/rc.d/rc.subr

pexp="nginx:.*master.*process.*$daemon"

rc_cmd $1
