# sysutils/monit/pkg/monit.rc
#!/bin/sh
#
# $OpenBSD: monit.rc,v 1.3 2014/09/24 09:30:48 jasper Exp $

daemon="${TRUEPREFIX}/bin/monit"

. /etc/rc.d/rc.subr

rc_cmd $1
