# sysutils/monit/pkg/monit.rc
#!/bin/sh
#
# $OpenBSD: monit.rc,v 1.2 2012/11/25 13:34:04 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/monit"
daemon_flags="-c ${SYSCONFDIR}/monitrc"

. /etc/rc.d/rc.subr

rc_cmd $1
