# sysutils/monit/pkg/monit.rc
#!/bin/ksh
#
# $OpenBSD: monit.rc,v 1.4 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/bin/monit"

. /etc/rc.d/rc.subr

rc_cmd $1
