# sysutils/monit/pkg/monit.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/monit"

. /etc/rc.d/rc.subr

rc_cmd $1
