# net/vnstat/pkg/vnstatd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/vnstatd -d"
daemon_user="_vnstat"

. /etc/rc.d/rc.subr

rc_pre() {
	/usr/bin/install -d -o _vnstat /var/run/vnstat
}

rc_cmd $1
