# net/syncthing/pkg/syncthing.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/syncthing"
daemon_flags="-no-browser"
daemon_user="_syncthing"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
