# net/syncthing/pkg/syncthing.rc
#!/bin/sh
#
# $OpenBSD: syncthing.rc,v 1.2 2016/04/23 08:02:06 edd Exp $

daemon="${TRUEPREFIX}/bin/syncthing"
daemon_flags="-no-browser"
daemon_user="_syncthing"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_bg=YES

rc_cmd $1
