# net/openvpn/pkg/openvpn.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/openvpn --daemon"

. /etc/rc.d/rc.subr

rc_cmd $1
