# net/pfstat/pkg/pfstatd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/pfstatd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
