# games/yquake2/pkg/q2ded.rc
#!/bin/sh
#
# $OpenBSD: q2ded.rc,v 1.2 2017/05/27 13:53:06 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/q2ded"
daemon_user="_yquake2"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_reload=NO

rc_cmd $1
