# net/minio/server/pkg/minio.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/minio server"
daemon_flags="${VARBASE}/minio/export"
daemon_user="_minio"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
