# net/endlessh/pkg/endlessh.rc
#!/bin/ksh

daemon="${LOCALBASE}/bin/endlessh"
daemon_flags="-s -4"
daemon_user=_endlessh
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
