# net/vsftpd/pkg/vsftpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/vsftpd -obackground=yes"
daemon_flags="${SYSCONFDIR}/vsftpd.conf"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
