# net/vsftpd/pkg/vsftpd.rc
#!/bin/sh
#
# $OpenBSD: vsftpd.rc,v 1.2 2014/01/30 14:01:50 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/vsftpd -obackground=yes"
daemon_flags="${SYSCONFDIR}/vsftpd.conf"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
