# sysutils/netshot/pkg/netshot.rc
#!/bin/ksh
#
# $OpenBSD: netshot.rc,v 1.1.1.1 2018/08/28 11:25:53 remi Exp $

daemon=$(${LOCALBASE}/bin/javaPathHelper -c netshot)
daemon_flags="-jar ${TRUEPREFIX}/netshot/netshot.jar"
daemon_user="_netshot"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
