# telephony/asterisk/16/pkg/asterisk.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/asterisk"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_start() {
	${rcexec} "${TRUEPREFIX}/sbin/safe_asterisk"
}

rc_stop() {
	${daemon} -rx "core stop now"
}

rc_cmd $1
