# devel/py-buildslave/pkg/buildslave.rc
#!/bin/sh
#
# $OpenBSD: buildslave.rc,v 1.3 2016/06/03 07:59:22 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/buildslave"
daemon_flags="${HOMEDIR}"
daemon_user="_buildslave"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon} start${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_start() {
	${rcexec} "${daemon} start ${daemon_flags}"
}

rc_stop() {
	${daemon} stop ${daemon_flags}
}

rc_cmd $1
