# devel/py-buildslave/pkg/buildslave.rc
#!/bin/ksh
#
# $OpenBSD: buildslave.rc,v 1.5 2018/01/11 19:27:02 rpe Exp $

daemon="${TRUEPREFIX}/bin/buildbot-worker"
daemon_flags="${HOMEDIR}"
daemon_user="_buildslave"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon} start${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_start() {
	${rcexec} "${daemon} start ${daemon_flags}"
}

rc_stop() {
	${daemon} stop ${daemon_flags}
}

rc_cmd $1
