# lang/php/8.5/pkg/php85_fpm.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/php-fpm-${PV}"

. /etc/rc.d/rc.subr

pexp="php-fpm-${PV}: master process .*"
rc_reload_signal=USR2

rc_cmd $1
