# net/munin/pkg/munin_asyncd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/munin/munin-asyncd --fork"
daemon_user="_munin"

. /etc/rc.d/rc.subr

pexp="perl: munin-asyncd .*"
rc_bg=YES

rc_cmd $1
