# net/munin/pkg/munin_node.rc
#!/bin/sh
#
# $OpenBSD: munin_node.rc,v 1.6 2013/01/08 11:14:02 kirby Exp $

daemon="${TRUEPREFIX}/sbin/munin-node"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl -wT $daemon"

rc_pre() {
        install -d -o _munin /var/run/munin
}

rc_cmd $1
