# net/munin/pkg/munin_asyncd.rc
#!/bin/ksh
#
# $OpenBSD: munin_asyncd.rc,v 1.5 2018/01/11 19:27:06 rpe Exp $

daemon="${TRUEPREFIX}/libexec/munin/munin-asyncd --fork"
daemon_user="_munin"

. /etc/rc.d/rc.subr

pexp="perl: munin-asyncd .*"
rc_bg=YES

rc_cmd $1
