# sysutils/salt/pkg/salt_syndic.rc
#!/bin/ksh
#
# $OpenBSD: salt_syndic.rc,v 1.3 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/bin/salt-syndic -d"
daemon_flags="-u _salt"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
