# sysutils/salt/pkg/salt_minion.rc
#!/bin/ksh
#
# $OpenBSD: salt_minion.rc,v 1.4 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/bin/salt-minion -d"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
