# sysutils/salt/pkg/salt_minion.rc
#!/bin/sh
#
# $OpenBSD: salt_minion.rc,v 1.3 2014/06/17 12:30:39 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/salt-minion -d"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
