# sysutils/salt/pkg/salt_minion.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/salt-minion -d"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
