# comms/hylafax/pkg/hfaxd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/hfaxd"
daemon_flags="-l 0.0.0.0 -i hylafax -l ::0 -i hylafax"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
