# mail/dspam/pkg/dspam.rc
#!/bin/ksh
#
# $OpenBSD: dspam.rc,v 1.2 2018/01/11 19:27:03 rpe Exp $

daemon="${TRUEPREFIX}/bin/dspam --daemon"
daemon_user="_dspam"

. /etc/rc.d/rc.subr

rc_cmd $1
