# www/honk/pkg/honk.rc
#!/bin/ksh
#
# $OpenBSD: honk.rc,v 1.5 2021/02/27 09:28:51 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/honk"
daemon_flags="-datadir ${VARBASE}/honk -viewdir ${VARBASE}/honk"
daemon_logger="daemon.info"
daemon_user="_honk"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_stop() {
	pkill -T "${daemon_rtable}" -U "${daemon_user}" -x honk
}

rc_cmd $1
