# www/honk/pkg/honk.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/honk"
daemon_flags="-datadir ${VARBASE}/honk -viewdir ${VARBASE}/honk"
daemon_logger="daemon.info"
daemon_user="_honk"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/bin/honk .*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
