# security/vault/pkg/vault_agent.rc
#!/bin/ksh

daemon="/usr/local/bin/vault agent"
daemon_flags="-config=${SYSCONFDIR}/vault/agent.hcl"
daemon_user="_vault"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_cmd $1
