# sysutils/riemann/pkg/riemann.rc
#!/bin/ksh
#
# $OpenBSD: riemann.rc,v 1.5 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/bin/riemann"
daemon_user="_riemann"

. /etc/rc.d/rc.subr

java="$(${LOCALBASE}/bin/javaPathHelper -c clojure)"

pexp="${java}.*riemann.*"
rc_bg=YES

rc_start() {
	${rcexec} "JAVA_CMD=\"${java}\" ${daemon}"
}

rc_cmd $1
