# sysutils/riemann/pkg/riemann.rc
#!/bin/sh
#
# $OpenBSD: riemann.rc,v 1.4 2017/05/27 19:57:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/riemann"
daemon_user="_riemann"

. /etc/rc.d/rc.subr

java="$(${LOCALBASE}/bin/javaPathHelper -c clojure)"

pexp="${java}.*riemann.*"
rc_bg=YES

rc_start() {
	${rcexec} "JAVA_CMD=\"${java}\" ${daemon}"
}

rc_cmd $1
