# x11/gnustep/base/pkg/gdomap.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/gdomap"
daemon_flags="-p"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
