# net/gdnsd/pkg/gdnsd.rc
#!/bin/sh
#
# $OpenBSD: gdnsd.rc,v 1.4 2017/08/10 20:14:56 jung Exp $

daemon="${PREFIX}/sbin/gdnsd"
daemon_flags="start"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_restart() {
	${daemon} restart
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
