# net/gdnsd/pkg/gdnsd.rc
#!/bin/ksh
#
# $OpenBSD: gdnsd.rc,v 1.5 2018/01/11 19:27:05 rpe Exp $

daemon="${PREFIX}/sbin/gdnsd"
daemon_flags="start"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_restart() {
	${daemon} restart
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
