# sysutils/syslog-ng/pkg/syslog_ng.rc
#!/bin/ksh
#
# $OpenBSD: syslog_ng.rc,v 1.8 2018/09/22 23:12:09 sthen Exp $

daemon="${TRUEPREFIX}/sbin/syslog-ng"

. /etc/rc.d/rc.subr

pexp="${daemon}"

rc_cmd $1
