# sysutils/syslog-ng/pkg/syslog_ng.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/syslog-ng"

. /etc/rc.d/rc.subr

pexp="${daemon}"

rc_cmd $1
