# sysutils/syslog-ng/pkg/syslog_ng.rc
#!/bin/sh
#
# $OpenBSD: syslog_ng.rc,v 1.6 2013/11/27 07:52:54 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/syslog-ng"
daemon_flags="-p /var/run/syslog.pid"

. /etc/rc.d/rc.subr

pexp="${daemon}"

rc_cmd $1
