# audio/zeya/pkg/zeya.rc
#!/bin/sh
#
# $OpenBSD: zeya.rc,v 1.3 2017/05/27 19:57:06 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/zeya"
daemon_flags="--path=/var/tmp"
daemon_user="_zeya"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_bg=YES
rc_reload=NO

rc_cmd $1
