# www/lighttpd/pkg/lighttpd.rc
#!/bin/ksh
#
# $OpenBSD: lighttpd.rc,v 1.4 2018/01/11 19:27:11 rpe Exp $

daemon="${TRUEPREFIX}/sbin/lighttpd"
daemon_flags="-f ${SYSCONFDIR}/lighttpd.conf"

. /etc/rc.d/rc.subr

rc_reload() {
	pkill -INT -xf "${pexp}"
	${rcexec} "${daemon} ${daemon_flags}" >/dev/null
}

rc_cmd $1
