# www/lighttpd/pkg/lighttpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/lighttpd"
daemon_flags="-f ${SYSCONFDIR}/lighttpd.conf"

. /etc/rc.d/rc.subr

rc_reload_signal=SIGUSR1
rc_stop_signal=INT

rc_cmd $1
