# www/lighttpd/pkg/lighttpd.rc
#!/bin/ksh
#
# $OpenBSD: lighttpd.rc,v 1.5 2021/10/31 17:38:51 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/lighttpd"
daemon_flags="-f ${SYSCONFDIR}/lighttpd.conf"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
