# net/iperf3/pkg/iperf3.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/iperf3 -s -D"
daemon_user="_iperf3"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
