# lang/php/5.5/pkg/php_fpm.rc
#!/bin/sh
#
# $OpenBSD: php_fpm.rc,v 1.2 2015/07/15 16:03:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/php-fpm-${PV}"

. /etc/rc.d/rc.subr

pexp="php-fpm-${PV}: master process .*"
rc_reload=NO

rc_cmd $1
