# www/sogo/pkg/sogod.rc
#!/bin/sh
#
# $OpenBSD: sogod.rc,v 1.2 2010/12/11 14:21:08 sebastia Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/sogod"
daemon_user="_sogo"

rc_cmd $1
