# www/sogo/pkg/sogod.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/sogod"
daemon_user="_sogo"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o _sogo -g _sogo /var/run/sogo/
}

rc_cmd $1
