# www/sogo/pkg/sogod.rc
#!/bin/sh
#
# $OpenBSD: sogod.rc,v 1.5 2011/04/18 12:48:29 sebastia Exp $

daemon="${TRUEPREFIX}/sbin/sogod"
daemon_user="_sogo"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
