# net/ladvd/pkg/ladvd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/ladvd"
daemon_flags="-a -z"

. /etc/rc.d/rc.subr

pexp="ladvd: parent \[priv\].*"
rc_reload=NO

rc_cmd $1
