# www/nginx/pkg/nginx.rc
#!/bin/sh
#
# $OpenBSD: nginx.rc,v 1.7 2011/03/09 19:37:47 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/nginx"

. /etc/rc.d/rc.subr

pexp="nginx:.*master.*process"

rc_reload() {
	${daemon} -s reload
}

rc_stop() {
	${daemon} -s stop
	[ $? -gt 0 ] && pkill -f "^${pexp}"
}

rc_cmd $1
