# www/nginx/pkg/nginx.rc
#!/bin/sh
#
# $OpenBSD: nginx.rc,v 1.12 2015/07/29 17:21:11 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/nginx"

. /etc/rc.d/rc.subr

pexp="nginx: master process ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
