# www/nginx/pkg/nginx.rc
#!/bin/sh
#
# $OpenBSD: nginx.rc,v 1.8 2011/04/09 01:21:55 william Exp $

daemon="${TRUEPREFIX}/sbin/nginx"

. /etc/rc.d/rc.subr

pexp="nginx:.*master.*process"

rc_reload() {
	${daemon} -s reload
}

rc_stop() {
	${daemon} -s stop || pkill -f "^${pexp}"
}

rc_cmd $1
