# net/icecast/pkg/icecast.rc
#!/bin/sh
#
# $OpenBSD: icecast.rc,v 1.2 2010/11/11 10:38:47 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/icecast"
daemon_flags="-b -c /var/icecast/icecast.xml"

rc_cmd $1
