# net/icecast/pkg/icecast.rc
#!/bin/sh
#
# $OpenBSD: icecast.rc,v 1.3 2010/12/24 10:40:05 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/icecast"
daemon_flags="-b -c /var/icecast/icecast.xml"

rc_conf

rc_cmd $1
