# sysutils/beats/hwsensorsbeat/pkg/hwsensorsbeat.rc
#!/bin/sh
#
# $OpenBSD: hwsensorsbeat.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${PREFIX}/bin/hwsensorsbeat"
daemon_flags="-c ${SYSCONFDIR}/hwsensorsbeat/hwsensorsbeat.yml"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
