# sysutils/obsdfreqd/pkg/obsdfreqd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/obsdfreqd"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_bg=YES

rc_cmd $1
