# net/nagios/nsca/pkg/nsca.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/nsca"
daemon_flags="-c ${SYSCONFDIR}/nsca.cfg"

. /etc/rc.d/rc.subr

rc_cmd $1
