# net/litecoin/pkg/litecoind.rc
#!/bin/ksh
#
# $OpenBSD: litecoind.rc,v 1.2 2019/07/07 15:22:20 rsadowski Exp $

daemon="${TRUEPREFIX}/bin/litecoind"
daemon_timeout="300"
daemon_user="_litecoind"
data_dir="${VARBASE}/litecoin"
daemon_flags="-daemon -conf=/etc/litecoin.conf -datadir=${data_dir}"

. /etc/rc.d/rc.subr

rc_reload=NO


rc_cmd $1
