# sysutils/slant/pkg/slant_collectd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/slant-collectd"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
