# sysutils/openpoppassd/pkg/openpoppassd.rc
#!/bin/ksh
#
# $OpenBSD: openpoppassd.rc,v 1.2 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/libexec/openpoppassd"

. /etc/rc.d/rc.subr

pexp="openpoppassd: \[master\]"
rc_reload=NO

rc_cmd $1
