# mail/amavisd-new/pkg/amavis_mc.rc
#!/bin/ksh
#
# $OpenBSD: amavis_mc.rc,v 1.3 2018/01/11 19:27:02 rpe Exp $

daemon="${TRUEPREFIX}/sbin/amavis-mc"

. /etc/rc.d/rc.subr

rc_reload=NO

pexp="/usr/bin/perl -T ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
