# security/vpnc/pkg/vpnc.rc
#!/bin/ksh
#
# $OpenBSD: vpnc.rc,v 1.2 2018/01/11 19:27:09 rpe Exp $

daemon="${TRUEPREFIX}/sbin/vpnc"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
