# net/spectrum/pkg/spectrum.rc
#!/bin/sh
#
# $OpenBSD: spectrum.rc,v 1.2 2011/03/07 10:14:01 sebastia Exp $

daemon="${TRUEPREFIX}/sbin/spectrumctl"

. /etc/rc.d/rc.subr

rc_start() {
	${daemon} start
}
rc_stop() {
	${daemon} stop
}
rc_reload() {
	${daemon} reload
}

rc_cmd $1
