# sysutils/smartmontools/pkg/smartd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/smartd"

. /etc/rc.d/rc.subr

rc_cmd $1
