# sysutils/smartmontools/pkg/smartd.rc
#!/bin/sh
#
# $OpenBSD: smartd.rc,v 1.2 2010/11/23 18:53:46 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/smartd"

rc_cmd $1
