# sysutils/smartmontools/pkg/smartd.rc
#!/bin/ksh
#
# $OpenBSD: smartd.rc,v 1.5 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/sbin/smartd"

. /etc/rc.d/rc.subr

rc_cmd $1
