# databases/mysql/pkg/mysqld.rc
#!/bin/sh
#
# $OpenBSD: mysqld.rc,v 1.11 2011/03/06 16:54:38 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mysqld_safe"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_check() {
	${TRUEPREFIX}/bin/mysqladmin -s ping >/dev/null
}

rc_stop() {
	${TRUEPREFIX}/bin/mysqladmin -s shutdown
	if [ $? -gt 0 ]; then
		pkill -f "^/bin/sh ${daemon}"
		pkill -U _mysql -f "^${TRUEPREFIX}/libexec/mysqld"
	fi
}

rc_cmd $1
