# telephony/asterisk/pkg/asterisk.rc
#!/bin/sh
#
# $OpenBSD: asterisk.rc,v 1.2 2011/03/09 18:51:49 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/asterisk"

. /etc/rc.d/rc.subr

rc_start() {
	${TRUEPREFIX}/sbin/safe_asterisk
}

rc_stop() {
	${daemon} -rx "core stop now"
}

rc_cmd $1
