# security/hitch/pkg/hitch.rc
#!/bin/ksh
#
# $OpenBSD: hitch.rc,v 1.2 2018/12/09 12:05:25 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/hitch --daemon"
daemon_flags="-u _hitch -g _hitch -s -r ${LOCALSTATEDIR}/hitch --config=${SYSCONFDIR}/hitch.conf"

. /etc/rc.d/rc.subr

rc_cmd $1
