# www/smarc/pkg/msearchd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/msearchd"

. /etc/rc.d/rc.subr

rc_cmd $1
