# mail/pop3d/pkg/pop3d.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/pop3d"

. /etc/rc.d/rc.subr

pexp="pop3d: \[priv\]"

rc_reload=NO

rc_cmd $1
