# mail/pop3d/pkg/pop3d.rc
#!/bin/ksh
#
# $OpenBSD: pop3d.rc,v 1.2 2018/01/11 19:27:03 rpe Exp $

daemon="${TRUEPREFIX}/sbin/pop3d"

. /etc/rc.d/rc.subr

pexp="pop3d: \[priv\]"

rc_reload=NO

rc_cmd $1
