# security/sshguard/pkg/sshguard.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/sshguard"

. /etc/rc.d/rc.subr

pexp="/bin/sh $daemon${daemon_flags:+ $daemon_flags}"

rc_bg=YES
rc_reload=NO

rc_cmd $1
