# security/sshguard/pkg/sshguard.rc
#!/bin/ksh
#
# $OpenBSD: sshguard.rc,v 1.5 2019/03/29 15:54:24 sthen Exp $

daemon="${TRUEPREFIX}/sbin/sshguard"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/libexec/sshg-blocker .*"

rc_bg=YES
rc_reload=NO

rc_cmd $1
