# security/sshguard/pkg/sshguard.rc
#!/bin/ksh
#
# $OpenBSD: sshguard.rc,v 1.6 2020/08/18 05:06:29 rsadowski Exp $

daemon="${TRUEPREFIX}/sbin/sshguard"

. /etc/rc.d/rc.subr

pexp="/bin/sh $daemon"

rc_bg=YES
rc_reload=NO

rc_cmd $1
