# net/totd/pkg/totd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/totd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
