# net/gelatod/pkg/gelatod.rc
#!/bin/ksh
#

daemon="/usr/local/sbin/gelatod"

. /etc/rc.d/rc.subr

rc_cmd $1
