# sysutils/restic-rest-server/pkg/restic_rest_server.rc
#!/bin/ksh
#
# $OpenBSD: restic_rest_server.rc,v 1.2 2019/03/05 13:49:50 abieber Exp $

daemon="${TRUEPREFIX}/bin/restic-rest-server"
daemon_flags="--path ${VARBASE}/restic"
daemon_user="_restic"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_cmd $1
