# security/sudo/pkg/sudo_logsrvd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/sudo_logsrvd"

. /etc/rc.d/rc.subr

rc_cmd $1
