# net/netatalk3/pkg/netatalk.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/netatalk"

. /etc/rc.d/rc.subr

rc_cmd $1
