# www/cherokee/pkg/cherokee.rc
#!/bin/sh
#
# $OpenBSD: cherokee.rc,v 1.2 2010/12/24 10:40:06 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/cherokee"
daemon_flags="-d"

rc_conf

pexp="${daemon}[^-]"

rc_cmd $1
