# www/cherokee/pkg/cherokee.rc
#!/bin/sh
#
# $OpenBSD: cherokee.rc,v 1.4 2011/12/17 22:01:13 sthen Exp $

daemon="${TRUEPREFIX}/sbin/cherokee -d"

. /etc/rc.d/rc.subr

pexp="${daemon}[^-]"

rc_cmd $1
