# net/ziproxy/pkg/ziproxy.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/ziproxy"
daemon_flags="-d > /var/run/ziproxy.pid"

. /etc/rc.d/rc.subr

rc_cmd $1
