# net/tor/pkg/tor.rc
#!/bin/sh
#
# $OpenBSD: tor.rc,v 1.1 2011/03/11 06:49:55 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/tor"

. /etc/rc.d/rc.subr

rc_cmd $1
