# net/tor/pkg/tor.rc
#!/bin/sh
#
# $OpenBSD: tor.rc,v 1.2 2015/04/03 15:26:10 pascal Exp $

daemon="${TRUEPREFIX}/bin/tor"
daemon_timeout=60

. /etc/rc.d/rc.subr

rc_stop() {
	pkill -INT -f "^${pexp}"
}

rc_cmd $1
