# net/tor/pkg/tor.rc
#!/bin/ksh
#
# $OpenBSD: tor.rc,v 1.7 2021/11/06 10:38:40 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/tor"
daemon_timeout=60

. /etc/rc.d/rc.subr

rc_stop_signal=INT

rc_cmd $1
