# net/tor/pkg/tor.rc
#!/bin/sh
#
# $OpenBSD: tor.rc,v 1.4 2015/07/15 13:54:55 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/tor"
daemon_timeout=60

. /etc/rc.d/rc.subr

rc_stop() {
	pkill -INT -xf "^${pexp}"
}

rc_cmd $1
