# sysutils/apcupsd/pkg/apcupsd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/apcupsd"
daemon_flags="--kill-on-powerfail"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	rm -f ${SYSCONFDIR}/powerfail
}

rc_cmd $1
