# games/hedgewars/pkg/hedgewars_server.rc
#!/bin/sh
#
# $OpenBSD: hedgewars_server.rc,v 1.2 2016/11/04 11:15:14 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/hedgewars-server"
daemon_user="_hedgewars"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
