# mail/cyrus-imapd/pkg/cyrus_imapd.rc
#!/bin/ksh
#
# $OpenBSD: cyrus_imapd.rc,v 1.13 2020/05/14 12:26:39 ajacoutot Exp $

daemon="${TRUEPREFIX}/cyrus/libexec/master -d"

. /etc/rc.d/rc.subr

rc_pre() {
	install -d /var/run/cyrus/socket/
	rm -f /var/run/cyrus-master.pid
}

rc_cmd $1
