# mail/cyrus-imapd/pkg/cyrus_imapd.rc
#!/bin/sh
#
# $OpenBSD: cyrus_imapd.rc,v 1.4 2010/12/04 14:11:53 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/libexec/cyrus-imapd/master"
daemon_flags="-d"

# needed when using the BerkeleyDB backend to prevent:
# DBERROR db4: no absolute path for the current directory: No such file or directory
rc_pre() {
	cd /var/imap
}

rc_post() {
	rm -f /var/run/cyrus-master.pid
}

rc_cmd $1
