# mail/cyrus-imapd/pkg/cyrus_imapd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/cyrus/libexec/master -d"

. /etc/rc.d/rc.subr

rc_pre() {
	install -d /var/run/cyrus/socket/
	rm -f /var/run/cyrus-master.pid
}

rc_cmd $1
