# mail/cyrus-imapd/pkg/cyrus_imapd.rc
#!/bin/ksh
#
# $OpenBSD: cyrus_imapd.rc,v 1.12 2018/01/11 19:27:03 rpe Exp $

daemon="${TRUEPREFIX}/cyrus/libexec/master -d"

. /etc/rc.d/rc.subr

rc_pre() {
	rm -f /var/run/cyrus-master.pid
}

rc_cmd $1
