# mail/cyrus-imapd/pkg/cyrus_imapd.rc
#!/bin/sh
#
# $OpenBSD: cyrus_imapd.rc,v 1.7 2010/12/27 14:50:23 ajacoutot Exp $

daemon="${TRUEPREFIX}/libexec/cyrus-imapd/master"
daemon_flags="-d"

. /etc/rc.d/rc.subr

# chdir needed when using the BerkeleyDB backend to prevent:
# DBERROR db4: no absolute path for the current directory: No such file or directory
rc_pre() {
	rm -f /var/run/cyrus-master.pid
	cd /var/imap
}

rc_cmd $1
