# sysutils/rsyslog/pkg/rsyslogd.rc
#!/bin/ksh
#
# $OpenBSD: rsyslogd.rc,v 1.3 2018/03/13 09:18:16 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/rsyslogd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
