# sysutils/rsyslog/pkg/rsyslogd.rc
#!/bin/ksh
#
# $OpenBSD: rsyslogd.rc,v 1.2 2018/03/13 06:05:34 remi Exp $

daemon="${TRUEPREFIX}/sbin/rsyslogd"

rc_reload=NO

. /etc/rc.d/rc.subr

rc_cmd $1
