# audio/mumble/pkg/murmurd.rc
#!/bin/ksh
#
# $OpenBSD: murmurd.rc,v 1.3 2022/02/18 12:59:18 kn Exp $

daemon="${TRUEPREFIX}/bin/murmurd"
daemon_flags="-ini ${SYSCONFDIR}/murmur.ini"
daemon_user="_murmurd"

. /etc/rc.d/rc.subr

# only reloads TLS settings
rc_reload_signal=USR1

rc_cmd $1
