# audio/mumble/pkg/murmurd.rc
#!/bin/ksh
#
# $OpenBSD: murmurd.rc,v 1.2 2021/11/06 10:38:41 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/murmurd"
daemon_flags="-ini ${SYSCONFDIR}/murmur.ini"
daemon_user="_murmurd"

. /etc/rc.d/rc.subr

# only reloads TLS settings
rc_reload_signal=USR1

rc_cmd $1
