# sysutils/serf/pkg/serf.rc
#!/bin/sh
#
# $OpenBSD: serf.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/serf agent"
daemon_flags="-syslog"
daemon_user="_serf"

. /etc/rc.d/rc.subr

rc_stop() {
	pkill -INT -T "${daemon_rtable}" -xf "${pexp}"
}

rc_cmd $1
