# databases/riak/pkg/riak.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/riak"
daemon_user="_riak"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_usercheck=NO

rc_check() {
	${rcexec} "${daemon} ping"
}

rc_start() {
	${rcexec} "${daemon} start"
}

rc_stop() {
	${rcexec} "${daemon} stop"
}

rc_cmd $1
