# net/netatalk/pkg/afpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/afpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
