# net/netatalk/pkg/cnid_metad.rc
#!/bin/sh
#
# $OpenBSD: cnid_metad.rc,v 1.1 2012/04/11 18:39:13 sthen Exp $

daemon="${TRUEPREFIX}/sbin/cnid_metad"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
