# net/netatalk/pkg/cnid_metad.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/cnid_metad"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
