# devel/snare/pkg/snare.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/snare"
daemon_user="_snare"

. /etc/rc.d/rc.subr

rc_cmd $1
