# net/dnsmasq/pkg/dnsmasq.rc
#!/bin/ksh
#
# $OpenBSD: dnsmasq.rc,v 1.2 2018/01/11 19:27:05 rpe Exp $

daemon="${TRUEPREFIX}/sbin/dnsmasq"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
