# news/sabnzbd/pkg/sabnzbd.rc
#!/bin/ksh

daemon="${MODPY_BIN} ${TRUEPREFIX}/lib/python${MODPY_VERSION}/site-packages/sabnzbd/SABnzbd.py -d"
daemon_flags="-b 0 -f ${LOCALSTATEDIR}/sabnzbd"
daemon_timeout=120
daemon_user="_sabnzbd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
