# news/sabnzbd/pkg/sabnzbd.rc
#!/bin/sh
#
# $OpenBSD: sabnzbd.rc,v 1.2 2015/07/16 00:12:02 ajacoutot Exp $

datadir="/var/sabnzbd"

daemon="${MODPY_BIN} ${TRUEPREFIX}/lib/python${MODPY_VERSION}/site-packages/sabnzbd/SABnzbd.py -d"
daemon_flags="-b 0 -f ${datadir} -p ${datadir}/sabnzbd.pid"
daemon_user="_sabnzbd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
