# net/smokeping/pkg/smokeping.rc
#!/bin/sh
#
# $OpenBSD: smokeping.rc,v 1.3 2014/11/07 23:56:23 sthen Exp $

daemon="${TRUEPREFIX}/bin/smokeping"
daemon_user="_smokeping"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl ${daemon}${daemon_flags:+ ${daemon_flags}}\$"

rc_pre() {
	install -d -o _smokeping /var/run/smokeping
}

rc_cmd $1
