# net/mattermost-server/pkg/mattermost.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/mattermost"
daemon_user="_mattermost"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
