# databases/mariadb/pkg/mysqld.rc
#!/bin/sh
#
# $OpenBSD: mysqld.rc,v 1.3 2017/05/27 13:53:06 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mysqld_safe"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/libexec/mysqld .*"
rc_reload=NO

rc_cmd $1
