# databases/mariadb/pkg/mysqld.rc
#!/bin/ksh
#
# $OpenBSD: mysqld.rc,v 1.5 2018/01/11 19:27:01 rpe Exp $

daemon="${TRUEPREFIX}/bin/mysqld_safe"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/libexec/mysqld .*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
