# databases/mariadb/pkg/mysqld.rc
#!/bin/ksh
#
# $OpenBSD: mysqld.rc,v 1.6 2020/07/07 12:40:45 sthen Exp $

daemon="${TRUEPREFIX}/bin/mysqld_safe"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/libexec/mariadbd .*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
