# databases/mariadb/pkg/mysqld.rc
#!/bin/sh
#
# $OpenBSD: mysqld.rc,v 1.2 2015/07/15 15:52:29 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mysqld_safe"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/libexec/mysqld .*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
