# mail/milter-checkrcpt/pkg/milter_checkrcpt.rc
#!/bin/ksh
#
# $OpenBSD: milter_checkrcpt.rc,v 1.2 2018/01/11 19:27:03 rpe Exp $

daemon="${TRUEPREFIX}/libexec/milter-checkrcpt"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
