# mail/milter-checkrcpt/pkg/milter_checkrcpt.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/milter-checkrcpt"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
