# net/zabbix/pkg/zabbix_proxy.rc
#!/bin/sh
#
# $OpenBSD: zabbix_proxy.rc,v 1.5 2010/12/27 14:50:23 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/zabbix_proxy"

. /etc/rc.d/rc.subr

pexp="zabbix_proxy: heartbeat sender"
rc_reload=NO

rc_pre() {
	install -d -o _zabbix -g _zabbix /var/run/zabbix
}

rc_cmd $1
