# net/zabbix/pkg/zabbix_proxy.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/zabbix_proxy"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o _zabbix -g _zabbix /var/run/zabbix
}

rc_cmd $1
