# net/zabbix/pkg/zabbix_server.rc
#!/bin/sh
#
# $OpenBSD: zabbix_server.rc,v 1.3 2010/12/11 12:57:37 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/zabbix_server"
pexp="zabbix_server: main process"
rc_reload=NO

rc_pre() {
	install -d -o _zabbix -g _zabbix /var/run/zabbix
}

rc_cmd $1
