# security/pcsc-lite/pkg/pcscd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/pcscd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
