# net/ntp/pkg/ntpd.rc
#!/bin/sh
#
# $OpenBSD: ntpd.rc,v 1.5 2010/12/27 14:50:23 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/ntpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
