# net/ntp/pkg/ntpd.rc
#!/bin/sh
#
# $OpenBSD: ntpd.rc,v 1.4 2010/12/24 10:40:05 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/ntpd"

rc_conf

rc_reload=NO

rc_cmd $1
