# net/softflowd/pkg/softflowd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/softflowd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
