# net/softflowd/pkg/softflowd.rc
#!/bin/ksh
#
# $OpenBSD: softflowd.rc,v 1.2 2018/01/11 19:27:08 rpe Exp $

daemon="${TRUEPREFIX}/sbin/softflowd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
