# sysutils/beats/metricbeat/pkg/metricbeat.rc
#!/bin/ksh
#
# $OpenBSD: metricbeat.rc,v 1.4 2018/01/11 19:27:09 rpe Exp $

daemon="${PREFIX}/bin/metricbeat"
daemon_flags="-path.config ${SYSCONFDIR}/metricbeat/"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
