# sysutils/beats/metricbeat/pkg/metricbeat.rc
#!/bin/sh
#
# $OpenBSD: metricbeat.rc,v 1.3 2017/05/27 19:57:07 ajacoutot Exp $

daemon="${PREFIX}/bin/metricbeat"
daemon_flags="-path.config ${SYSCONFDIR}/metricbeat/"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
