# net/tacacs+/pkg/tac_plus.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/tac_plus"
daemon_flags="-C ${SYSCONFDIR}/tac_plus.conf -Q _tacacs -U _tacacs"

. /etc/rc.d/rc.subr

rc_cmd $1
