# net/quagga/pkg/quagga_ripngd.rc
#!/bin/sh
#
# $OpenBSD: quagga_ripngd.rc,v 1.1 2013/04/18 20:29:08 sthen Exp $

daemon="${TRUEPREFIX}/sbin/ripngd -d"

. /etc/rc.d/rc.subr

rc_cmd $1
