# net/quagga/pkg/quagga_pimd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/pimd -d"

. /etc/rc.d/rc.subr

rc_cmd $1
