# net/quagga/pkg/quagga_ripd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/ripd -d"

. /etc/rc.d/rc.subr

rc_cmd $1
