# net/quagga/pkg/quagga_ospf6d.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/ospf6d -d"

. /etc/rc.d/rc.subr

rc_cmd $1
