# net/quagga/pkg/quagga_bgpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/bgpd -d"

. /etc/rc.d/rc.subr

rc_cmd $1
