# net/quagga/pkg/quagga_ripngd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/ripngd -d"

. /etc/rc.d/rc.subr

rc_cmd $1
