# net/quagga/pkg/quagga_babeld.rc
#!/bin/ksh
#
# $OpenBSD: quagga_babeld.rc,v 1.2 2018/01/11 19:27:07 rpe Exp $

daemon="${TRUEPREFIX}/sbin/babeld -d"

. /etc/rc.d/rc.subr

rc_cmd $1
