# net/quagga/pkg/quagga_ospf6d.rc
#!/bin/sh
#
# $OpenBSD: quagga_ospf6d.rc,v 1.1 2013/04/18 20:29:08 sthen Exp $

daemon="${TRUEPREFIX}/sbin/ospf6d -d"

. /etc/rc.d/rc.subr

rc_cmd $1
