# net/icbd/pkg/icbd.rc
#!/bin/sh
#
# $OpenBSD: icbd.rc,v 1.2 2015/07/15 16:22:45 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/icbd"

. /etc/rc.d/rc.subr

pexp="icbd: icbd"
rc_reload=NO

rc_cmd $1
