# net/icbd/pkg/icbd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/icbd"

. /etc/rc.d/rc.subr

pexp="icbd: icbd"
rc_reload=NO

rc_cmd $1
