# x11/gnome/gdm/pkg/gdm.rc
#!/bin/sh
#
# $OpenBSD: gdm.rc,v 1.3 2010/11/25 15:27:30 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/gdm"
pexp="${TRUEPREFIX}/sbin/gdm-binary"

rc_start() {
	(sleep 5; ${daemon})&
}

rc_reload() {
	rc_err "$0: reload is not supported."
}

rc_cmd $1
