# x11/gnome/gdm/pkg/gdm.rc
#!/bin/sh
#
# $OpenBSD: gdm.rc,v 1.4 2010/12/11 07:40:21 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/gdm"
pexp="${TRUEPREFIX}/sbin/gdm-binary"

rc_start() {
	(sleep 5; su -c ${_rcclass:-daemon} -s ${_rcshell:-/bin/sh} \
		${_rcuser:-root} -c "${daemon} ${daemon_flags}" >/dev/null)&
}

rc_reload() {
	rc_err "$0: reload is not supported."
}

rc_cmd $1
