# x11/gnome/gdm/pkg/gdm.rc
#!/bin/sh
#
# $OpenBSD: gdm.rc,v 1.21 2013/09/28 16:22:40 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/gdm"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -g _gdm -m 0711 /var/run/gdm
	install -d -o _gdm -g _gdm /var/run/gdm/greeter
}

rc_start() {
	${rcexec} "sleep 5; ${daemon} ${daemon_flags}" &
}

rc_stop() {
	pkill -f "^${pexp}"
	pkill -u _gdm
}

rc_cmd $1
