# x11/gnome/gdm/pkg/gdm.rc
#!/bin/sh
#
# $OpenBSD: gdm.rc,v 1.11 2011/03/09 18:51:49 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/gdm"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/sbin/gdm-binary"
rc_reload=NO

rc_start() {
	(sleep 5; ${rcexec} "${daemon} ${daemon_flags}")&
}

rc_cmd $1
