# x11/gnome/gdm/pkg/gdm.rc
#!/bin/sh
#
# $OpenBSD: gdm.rc,v 1.24 2013/12/01 17:54:13 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/gdm"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	( i=0
	while [ $i -le 10 ]; do
		pgrep -qf "^/usr/libexec/getty " && break
		sleep 1
		i=$((i+1))
	done
	[[ $i -gt 10 ]] && return 1
	${rcexec} "${daemon} ${daemon_flags}" ) &
}

rc_stop() {
	pkill -f "^${pexp}"
	pkill -u _gdm
}

rc_cmd $1
