# x11/gnome/gdm/pkg/gdm.rc
#!/bin/sh
#
# $OpenBSD: gdm.rc,v 1.19 2013/04/03 15:45:47 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/gdm"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_stop=NO

rc_pre() {
	install -d -g _gdm -m 0711 /var/run/gdm
	install -d -o _gdm -g _gdm /var/run/gdm/greeter
}

rc_start() {
	(sleep 5; ${rcexec} "${daemon} ${daemon_flags}")&
}

rc_cmd $1
