# x11/gnome/gdm/pkg/gdm.rc
#!/bin/sh
#
# $OpenBSD: gdm.rc,v 1.17 2012/10/16 16:38:13 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/gdm"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/sbin/gdm-binary"
rc_reload=NO

rc_pre() {
	install -d -g _gdm -m 0711 /var/run/gdm
	install -d -o _gdm -g _gdm /var/run/gdm/greeter
}

rc_post() {
	pkill -f "^gdm-session-worker"
}

rc_cmd $1
