# www/pnp4nagios/pkg/npcd.rc
#!/bin/sh
#
# $OpenBSD: npcd.rc,v 1.1 2016/05/18 09:24:00 sebastia Exp $

daemon="/usr/local/bin/npcd -d"
daemon_flags="-f /etc/pnp4nagios/npcd.cfg"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
