# www/pnp4nagios/pkg/npcd.rc
#!/bin/sh
#
# $OpenBSD: npcd.rc,v 1.2 2016/05/18 17:34:00 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/npcd -d"
daemon_flags="-f ${SYSCONFDIR}/pnp4nagios/npcd.cfg"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
