# www/pnp4nagios/pkg/npcd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/npcd -d"
daemon_flags="-f ${SYSCONFDIR}/npcd.cfg"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
