# audio/ympd/pkg/ympd.rc
#!/bin/sh
#
# $OpenBSD: ympd.rc,v 1.3 2017/05/27 13:53:06 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/ympd"
daemon_flags="-u _ympd"

. /etc/rc.d/rc.subr

rc_cmd $1
