# mail/zarafa/zarafa/pkg/zarafa_ical.rc
#!/bin/sh
#
# $OpenBSD: zarafa_ical.rc,v 1.2 2010/10/29 12:54:44 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/zarafa-ical"

rc_cmd $1
