# mail/zarafa/zarafa/pkg/zarafa-monitor.rc
#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/zarafa-monitor"

rc_cmd $1
