# mail/zarafa/zarafa/pkg/zarafa-ical.rc
#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/zarafa-ical"

rc_cmd $1
