# mail/zarafa/zarafa/pkg/zarafa_gateway.rc
#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/zarafa-gateway"

rc_cmd $1
