# mail/zarafa/zarafa/pkg/zarafa_monitor.rc
#!/bin/sh
#
# $OpenBSD: zarafa_monitor.rc,v 1.5 2015/11/02 11:06:34 robert Exp $

daemon="${TRUEPREFIX}/sbin/zarafa-monitor"

. /etc/rc.d/rc.subr

rc_cmd $1
