# sysutils/node_exporter/pkg/node_exporter.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/node_exporter"
daemon_logger="daemon.info"
daemon_user="_nodeexporter"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
