# www/wwwoffle/pkg/wwwoffled.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/wwwoffled"

. /etc/rc.d/rc.subr

rc_cmd $1
