# mail/popa3d/pkg/popa3d.rc
#!/bin/ksh

daemon="${PREFIX}/libexec/popa3d -D"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
