# mail/exim/pkg/exim.rc
#!/bin/sh
#
# $OpenBSD: exim.rc,v 1.2 2011/12/17 22:01:13 sthen Exp $

daemon="${TRUEPREFIX}/bin/exim -bd"
daemon_flags="-q30m"

. /etc/rc.d/rc.subr

pexp="${daemon}"

rc_pre() {
	[ ! -f /etc/mailer.conf.exim ]
}

rc_cmd $1
