# www/ruby-puma/pkg/PLIST
@comment $OpenBSD: PLIST,v 1.5 2019/11/12 06:48:14 jeremy Exp $
${GEM_BIN}/puma${GEM_BIN_SUFFIX}
${GEM_BIN}/pumactl${GEM_BIN_SUFFIX}
${GEM_LIB}/cache/${DISTNAME}.gem
${GEM_LIB}/gems/${DISTNAME}/
${GEM_LIB}/gems/${DISTNAME}/History.md
${GEM_LIB}/gems/${DISTNAME}/LICENSE
${GEM_LIB}/gems/${DISTNAME}/README.md
${GEM_LIB}/gems/${DISTNAME}/bin/
${GEM_LIB}/gems/${DISTNAME}/bin/puma
${GEM_LIB}/gems/${DISTNAME}/bin/puma-wild
${GEM_LIB}/gems/${DISTNAME}/bin/pumactl
${GEM_LIB}/gems/${DISTNAME}/docs/
${GEM_LIB}/gems/${DISTNAME}/docs/architecture.md
${GEM_LIB}/gems/${DISTNAME}/docs/deployment.md
${GEM_LIB}/gems/${DISTNAME}/docs/images/
${GEM_LIB}/gems/${DISTNAME}/docs/images/puma-connection-flow-no-reactor.png
${GEM_LIB}/gems/${DISTNAME}/docs/images/puma-connection-flow.png
${GEM_LIB}/gems/${DISTNAME}/docs/images/puma-general-arch.png
${GEM_LIB}/gems/${DISTNAME}/docs/nginx.md
${GEM_LIB}/gems/${DISTNAME}/docs/plugins.md
${GEM_LIB}/gems/${DISTNAME}/docs/restart.md
${GEM_LIB}/gems/${DISTNAME}/docs/signals.md
${GEM_LIB}/gems/${DISTNAME}/docs/systemd.md
${GEM_LIB}/gems/${DISTNAME}/docs/tcp_mode.md
${GEM_LIB}/gems/${DISTNAME}/lib/
${GEM_LIB}/gems/${DISTNAME}/lib/puma/
${GEM_LIB}/gems/${DISTNAME}/lib/puma.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/accept_nonblock.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/app/
${GEM_LIB}/gems/${DISTNAME}/lib/puma/app/status.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/binder.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/cli.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/client.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/cluster.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/commonlogger.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/configuration.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/const.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/control_cli.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/detect.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/dsl.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/events.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/io_buffer.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/jruby_restart.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/launcher.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/minissl/
${GEM_LIB}/gems/${DISTNAME}/lib/puma/minissl.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/minissl/context_builder.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/null_io.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/plugin/
${GEM_LIB}/gems/${DISTNAME}/lib/puma/plugin.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/plugin/tmp_restart.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/puma_http11.so
${GEM_LIB}/gems/${DISTNAME}/lib/puma/rack/
${GEM_LIB}/gems/${DISTNAME}/lib/puma/rack/builder.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/rack/urlmap.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/rack_default.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/reactor.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/runner.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/server.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/single.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/state_file.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/tcp_logger.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/thread_pool.rb
${GEM_LIB}/gems/${DISTNAME}/lib/puma/util.rb
${GEM_LIB}/gems/${DISTNAME}/lib/rack/
${GEM_LIB}/gems/${DISTNAME}/lib/rack/handler/
${GEM_LIB}/gems/${DISTNAME}/lib/rack/handler/puma.rb
${GEM_LIB}/gems/${DISTNAME}/tools/
${GEM_LIB}/gems/${DISTNAME}/tools/docker/
${GEM_LIB}/gems/${DISTNAME}/tools/docker/Dockerfile
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/README.md
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/init.d/
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/init.d/README.md
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/init.d/puma
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/init.d/run-puma
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/rc.d/
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/rc.d/README.md
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/rc.d/puma
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/rc.d/puma.conf
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/upstart/
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/upstart/README.md
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/upstart/puma-manager.conf
${GEM_LIB}/gems/${DISTNAME}/tools/jungle/upstart/puma.conf
${GEM_LIB}/gems/${DISTNAME}/tools/trickletest.rb
${GEM_LIB}/specifications/${DISTNAME}.gemspec
