# net/argus/pkg/argus.rc
#!/bin/sh
#
# $OpenBSD: argus.rc,v 1.5 2010/12/27 14:50:22 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/argus"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
