# net/argus/pkg/argus.rc
#!/bin/ksh
#
# $OpenBSD: argus.rc,v 1.6 2018/01/11 19:27:04 rpe Exp $

daemon="${TRUEPREFIX}/sbin/argus"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
