# net/argus/pkg/argus.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/argus"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
