# net/snort/pkg/snort.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/snort -D"
daemon_flags="-c ${SYSCONFDIR}/snort/snort.conf -u _snort -g _snort -t ${LOCALSTATEDIR}/snort -l ${LOCALSTATEDIR}/snort/log"
daemon_timeout=120

. /etc/rc.d/rc.subr

rc_cmd $1
