# devel/jenkins/pkg/jenkins.rc
#!/bin/sh
#
# $OpenBSD: jenkins.rc,v 1.4 2014/01/05 18:37:51 jasper Exp $

daemon=$(${LOCALBASE}/bin/javaPathHelper -c jenkins)
daemon_flags="-jar ${TRUEPREFIX}/share/jenkins/jenkins.war --httpPort=8000 --httpsPort=-1 --ajp13Port=-1"
daemon_user="_jenkins"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_stop {
	pkill -9 -f "^${pexp}"	
}

rc_cmd $1
