# databases/influxdb/pkg/influxdb.rc
#!/bin/ksh
#
# $OpenBSD: influxdb.rc,v 1.2 2018/01/11 19:27:01 rpe Exp $

daemon="${TRUEPREFIX}/bin/influxd"
daemon_user="_influx"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
