# sysutils/mcollective/pkg/mcollectived.rc
#!/bin/sh
#
# $OpenBSD: mcollectived.rc,v 1.2 2013/03/13 22:54:26 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/mcollectived"

. /etc/rc.d/rc.subr

pexp="${RUBY} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
