# sysutils/mcollective/pkg/mcollectived.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/mcollectived"
daemon_flags="--pid=/var/run/mcollectived.pid"

. /etc/rc.d/rc.subr

pexp="${RUBY} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
