# databases/redis/pkg/redis.rc
#!/bin/sh
#
# $OpenBSD: redis.rc,v 1.5 2015/07/15 16:15:35 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/redis-server"
daemon_flags="${SYSCONFDIR}/redis/redis.conf"
daemon_user="_redis"

. /etc/rc.d/rc.subr

pexp="redis-server: ${daemon} .*"

rc_reload=NO

rc_pre() {
	/usr/bin/install -d -o _redis -g _redis -m 0750 /var/run/redis
}

rc_cmd $1
