# print/cups-filters/pkg/cups_browsed.rc
#!/bin/sh
#
# $OpenBSD: cups_browsed.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/cups-browsed"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
