# print/cups-filters/pkg/cups_browsed.rc
#!/bin/ksh
#
# $OpenBSD: cups_browsed.rc,v 1.4 2018/01/11 19:27:08 rpe Exp $

daemon="${TRUEPREFIX}/sbin/cups-browsed"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
