# net/dico/pkg/dicod.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/dicod"
daemon_user="_dicod"

. /etc/rc.d/rc.subr

rc_cmd $1
