# mail/courier-authlib/pkg/courier_authdaemond.rc
#!/bin/sh
#
# $OpenBSD: courier_authdaemond.rc,v 1.4 2011/03/17 16:49:27 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/authdaemond"

. /etc/rc.d/rc.subr

pexp=".*courierlogger.*-start ${TRUEPREFIX}/libexec/courier-authlib/authdaemond"

rc_start() {
	mkdir -p /var/run/courier{,-auth}/
	${daemon} start
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
