# mail/courier-authlib/pkg/courier_authdaemond.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/authdaemond"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/sbin/courierlogger .* -start ${TRUEPREFIX}/libexec/courier-authlib/authdaemond"

rc_start() {
	mkdir -p /var/run/courier{,-auth}/
	${rcexec} "${daemon} start"
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
