# sysutils/dinit/pkg/dinit.rc
#!/bin/ksh
#
# $OpenBSD: dinit.rc,v 1.1.1.1 2019/06/07 15:01:52 edd Exp $

daemon="${TRUEPREFIX}/sbin/dinit"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_cmd $1
