# www/gotosocial/pkg/gotosocial.rc
#!/bin/ksh

daemon="${LOCALBASE}/bin/gotosocial"
daemon_user="_gotosocial"
daemon_flags="--config-path ${SYSCONFDIR}/gotosocial.yaml"
daemon_logger="daemon.info"

. /etc/rc.d/rc.subr

pexp="${daemon} .*"
rc_bg=YES
rc_reload=NO

rc_start() {
	rc_exec "${daemon} ${daemon_flags} server start"
}

rc_cmd $1
