# net/freeradius/pkg/freeradius.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/radiusd"

. /etc/rc.d/rc.subr

rc_pre() {
	/usr/bin/install -d -o _freeradius /var/run/radiusd
}

rc_cmd $1
