# net/mfi/pkg/mfi.rc
#!/bin/sh
#
# $OpenBSD: mfi.rc,v 1.4 2017/05/27 14:16:49 ajacoutot Exp $

daemon="${TRUEPREFIX}/share/mfi/lib/ace.jar"
daemon_user="_mfi"

. /etc/rc.d/rc.subr

java="$(${LOCALBASE}/bin/javaPathHelper -c mfi)"
pexp="${java} -jar ${daemon} start"

rc_reload=NO

rc_start() {
	${rcexec} "${java} ${daemon_flags} -jar ${daemon} start"
}

rc_stop() {
	${java} -jar ${daemon} stop
}

rc_cmd $1
