# net/bird/pkg/bird6.rc
#!/bin/sh
#
# $OpenBSD: bird6.rc,v 1.4 2014/07/01 21:51:42 sthen Exp $

daemon="${TRUEPREFIX}/sbin/bird6"

. /etc/rc.d/rc.subr

pexp="${daemon}${daemon_flags:+ ${daemon_flags}}$"

rc_stop() {
	${TRUEPREFIX}/sbin/birdc6 down
}

rc_cmd $1
