# net/bird/pkg/bird6.rc
#!/bin/sh
#
# $OpenBSD: bird6.rc,v 1.3 2014/07/01 20:38:44 sthen Exp $

daemon="${TRUEPREFIX}/sbin/bird6"

. /etc/rc.d/rc.subr

rc_check() {
        pgrep -q -f "^${pexp}$"
}

rc_stop() {
	${TRUEPREFIX}/sbin/birdc6 down
}

rc_cmd $1
