# net/bird/pkg/bird6.rc
#!/bin/ksh
#
# $OpenBSD: bird6.rc,v 1.6 2018/01/11 19:27:04 rpe Exp $

daemon="${TRUEPREFIX}/sbin/bird6"

. /etc/rc.d/rc.subr

rc_stop() {
	${TRUEPREFIX}/sbin/birdc6 down
}

rc_cmd $1
