# mail/mimedefang/pkg/mimedefang_multiplexor.rc
#!/bin/sh
#
# $OpenBSD: mimedefang_multiplexor.rc,v 1.9 2012/12/09 07:35:48 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mimedefang-multiplexor"
daemon_flags="-G -U _mdefang -E -l -m 2 -x 10 -b 600 -p /var/spool/MIMEDefang/mimedefang-multiplexor.pid"

. /etc/rc.d/rc.subr

rc_reload() {
	${TRUEPREFIX}/bin/md-mx-ctrl reread
}

rc_cmd $1
