# mail/mimedefang/pkg/mimedefang_multiplexor.rc
#!/bin/sh
#
# $OpenBSD: mimedefang_multiplexor.rc,v 1.5 2012/05/31 18:31:42 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mimedefang-multiplexor"
daemon_flags="-U _mdefang -l -m 2 -b 600 -p /var/spool/MIMEDefang/mimedefang-multiplexor.pid"

. /etc/rc.d/rc.subr

rc_reload() {
	pkill -INT -f "^${pexp}"
}

rc_cmd $1
