# mail/mimedefang/pkg/mimedefang.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/mimedefang"
daemon_flags="-Gqr -U _mdefang -p /var/spool/MIMEDefang/mimedefang.sock -m /var/spool/MIMEDefang/mimedefang-multiplexor.sock"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
