# mail/mimedefang/pkg/mimedefang_daemon.rc
#!/bin/sh
#
# $OpenBSD: mimedefang_daemon.rc,v 1.7 2012/09/18 15:34:39 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/mimedefang"
daemon_flags="-G -U _mdefang -P /var/spool/MIMEDefang/mimedefang.pid -p /var/spool/MIMEDefang/mimedefang.sock -m /var/spool/MIMEDefang/mimedefang-multiplexor.sock"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
