# mail/greyscanner/pkg/greyscanner.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/greyscanner"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl ${daemon}"
rc_reload=NO

rc_cmd $1
