# net/arpwatch/pkg/arpwatch.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/arpwatch"
daemon_flags="-u _arpwatch"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
