# net/arpwatch/pkg/arpwatch.rc
#!/bin/ksh
#
# $OpenBSD: arpwatch.rc,v 1.3 2018/10/31 08:58:56 sebastia Exp $

daemon="${TRUEPREFIX}/sbin/arpwatch"
daemon_flags="-u _arpwatch"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
