# comms/conserver/pkg/conserver.rc
#!/bin/sh
#
# $OpenBSD: conserver.rc,v 1.4 2010/12/27 14:50:22 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/conserver"
daemon_flags="-d -n"

. /etc/rc.d/rc.subr

rc_cmd $1
