# comms/conserver/pkg/conserver.rc
#!/bin/sh
#
# $OpenBSD: conserver.rc,v 1.5 2011/12/17 22:01:12 sthen Exp $

daemon="${TRUEPREFIX}/sbin/conserver -d"
daemon_flags="-n"

. /etc/rc.d/rc.subr

rc_cmd $1
