# comms/conserver/pkg/conserver.rc
#!/bin/ksh
#
# $OpenBSD: conserver.rc,v 1.9 2018/05/22 14:57:44 sthen Exp $

daemon="${TRUEPREFIX}/sbin/conserver -d"
daemon_flags="-n"
daemon_user="_conserver"

. /etc/rc.d/rc.subr

pexp="(${daemon}${daemon_flags:+ ${daemon_flags}}|conserver: master.*)"

rc_pre() {
	install -d -o _conserver /var/run/conserver
	[[ -f /var/log/conserver ]] || install -o _conserver \
	    -m 644 /dev/null /var/log/conserver
}

rc_cmd $1
