# comms/conserver/pkg/conserver.rc
#!/bin/sh
#
# $OpenBSD: conserver.rc,v 1.2 2010/12/11 08:05:38 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/conserver"
daemon_flags="-d -n"

rc_cmd $1
