# net/ircd-hybrid/pkg/ircd_hybrid.rc
#!/bin/sh
#
# $OpenBSD: ircd_hybrid.rc,v 1.1 2010/11/13 06:52:41 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/ircd"

rc_pre() {
	install -d -o _ircd-hybrid -g _ircd-hybrid /var/run/ircd
}

rc_cmd $1
