# multimedia/minidlna/pkg/minidlna.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/minidlnad"
daemon_user="_minidlna"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o ${daemon_user} /var/run/minidlna
}

rc_cmd $1
