# multimedia/minidlna/pkg/minidlna.rc
#!/bin/sh
#
# $OpenBSD: minidlna.rc,v 1.2 2013/05/07 20:37:24 sthen Exp $

daemon="${TRUEPREFIX}/sbin/minidlnad"
daemon_user="_minidlna"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o _minidlna /var/run/minidlna
}

rc_cmd $1
