# multimedia/minidlna/pkg/minidlna.rc
#!/bin/ksh
#
# $OpenBSD: minidlna.rc,v 1.3 2018/01/11 19:27:04 rpe Exp $

daemon="${TRUEPREFIX}/sbin/minidlnad"
daemon_user="_minidlna"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	install -d -o _minidlna /var/run/minidlna
}

rc_cmd $1
