# sysutils/ruby-puppet/pkg/puppetmasterd.rc
#!/bin/sh
#
# $OpenBSD: puppetmasterd.rc,v 1.2 2012/10/28 18:48:45 robert Exp $

daemon="/usr/local/bin/puppet master"

. /etc/rc.d/rc.subr

pexp=".*ruby.* ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
