# sysutils/ruby-puppet/pkg/puppet.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/puppet agent"

. /etc/rc.d/rc.subr

pexp="${RUBY} -W0 ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
