# net/pdnsd/pkg/pdnsd.rc
#!/bin/ksh
#
# $OpenBSD: pdnsd.rc,v 1.3 2018/01/11 19:27:07 rpe Exp $

daemon="${TRUEPREFIX}/sbin/pdnsd -d"

. /etc/rc.d/rc.subr

rc_reload() {
	${TRUEPREFIX}/sbin/pdnsd-ctl config
}

rc_cmd $1
