# mail/milter-regex/pkg/milter_regex.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/milter-regex"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
