# mail/milter-regex/pkg/milter_regex.rc
#!/bin/ksh
#
# $OpenBSD: milter_regex.rc,v 1.5 2018/01/11 19:27:03 rpe Exp $

daemon="${TRUEPREFIX}/libexec/milter-regex"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
