# security/clamav/pkg/clamav_milter.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/clamav-milter"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
