# security/clamav/pkg/clamd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/clamd"
daemon_timeout="75"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
