# security/clamav/pkg/clamd.rc
#!/bin/sh
#
# $OpenBSD: clamd.rc,v 1.2 2017/05/28 20:12:51 sthen Exp $

daemon="${TRUEPREFIX}/sbin/clamd"
daemon_timeout="90"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
