# security/clamav/pkg/clamav_freshclam.rc
#!/bin/sh
#
# $OpenBSD: clamav_freshclam.rc,v 1.1 2010/11/04 11:18:12 sthen Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/freshclam"
daemon_flags="-d"

rc_reload() {
	rc_err "$0: reload is not supported"
}

rc_cmd $1
