# security/clamav/pkg/clamd.rc
#!/bin/sh
#
# $OpenBSD: clamd.rc,v 1.1 2011/01/05 06:06:49 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/clamd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
