# security/clamav/pkg/clamav_daemon.rc
#!/bin/sh
#
# $OpenBSD: clamav_daemon.rc,v 1.4 2010/12/27 14:50:23 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/clamd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
