# security/clamav/pkg/clamd.rc
#!/bin/ksh
#
# $OpenBSD: clamd.rc,v 1.4 2019/10/20 10:57:14 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/clamd"
daemon_timeout="120"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
