# security/clamav/pkg/clamd.rc
#!/bin/ksh
#
# $OpenBSD: clamd.rc,v 1.5 2019/11/20 19:00:24 sthen Exp $

daemon="${TRUEPREFIX}/sbin/clamd"
daemon_timeout="75"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
