# sysutils/gpioflicker/pkg/gpioflicker.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/gpioflicker"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
