# sysutils/sysmon/pkg/sysmond.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/sysmond"
daemon_flags="-f ${SYSCONFDIR}/sysmon.conf"
daemon_timeout="60"

. /etc/rc.d/rc.subr

rc_reload() {
	${daemon} ${daemon_flags} -t && \
		pkill -HUP -T "${daemon_rtable}" -xf "${pexp}"
}

rc_cmd $1
