# sysutils/sysmon/pkg/sysmond.rc
#!/bin/ksh
#
# $OpenBSD: sysmond.rc,v 1.3 2018/01/11 19:27:11 rpe Exp $

daemon="${TRUEPREFIX}/bin/sysmond"
daemon_flags="-f ${SYSCONFDIR}/sysmon.conf"
daemon_timeout="60"

. /etc/rc.d/rc.subr

rc_reload() {
	${daemon} ${daemon_flags} -t && \
		pkill -HUP -T "${daemon_rtable}" -xf "${pexp}"
}

rc_cmd $1
