# misc/srcpd/pkg/srcpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/srcpd"

. /etc/rc.d/rc.subr

rc_stop_signal=KILL

rc_cmd $1
