# misc/srcpd/pkg/srcpd.rc
#!/bin/sh
#
# $OpenBSD: srcpd.rc,v 1.3 2010/12/24 10:40:06 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/srcpd"

rc_conf

rc_stop() {
	pkill -9 -f "^${pexp}"
}

rc_cmd $1
