# misc/srcpd/pkg/srcpd.rc
#!/bin/ksh
#
# $OpenBSD: srcpd.rc,v 1.7 2021/11/06 10:38:40 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/srcpd"

. /etc/rc.d/rc.subr

rc_stop_signal=KILL

rc_cmd $1
