# telephony/coturn/pkg/turnserver.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/turnserver --daemon"
daemon_user="_turnserver"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
