# net/unifi/pkg/unifi.rc
#!/bin/sh
#
# $OpenBSD: unifi.rc,v 1.3 2014/10/03 11:53:23 jasper Exp $

daemon="${TRUEPREFIX}/share/unifi/lib/ace.jar"
daemon_user="_unifi"

. /etc/rc.d/rc.subr

java="$(${LOCALBASE}/bin/javaPathHelper -c unifi)"
pexp="${java} -jar ${daemon} start"

rc_reload=NO

rc_start() {
	${rcexec} "${java} -jar ${daemon} start" &
}

rc_stop() {
	${java} -jar ${daemon} stop
}

rc_cmd $1
