# devel/got/pkg/gotd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/gotd"

. /etc/rc.d/rc.subr

pexp="gotd: parent"

rc_configtest() {
	${daemon} -n ${daemon_flags}
}

rc_reload=NO

rc_cmd $1
