# sysutils/loki/pkg/loki.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/loki"
daemon_user="_loki"
daemon_flags="-config.file ${SYSCONFDIR}/loki/loki-config.yaml"
daemon_logger="daemon.info"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
