# news/nzbget/pkg/nzbget.rc
#!/bin/ksh
daemon="${TRUEPREFIX}/bin/nzbget -D"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_stop() {
	${TRUEPREFIX}/bin/nzbget -Q
}

rc_cmd $1
