# x11/cde/pkg/dtlogin.rc
#!/bin/ksh
#

daemon="${TRUEPREFIX}/bin/dtlogin -daemon"

. /etc/rc.d/rc.subr

rc_cmd $1
