# games/ioquake3/pkg/ioq3ded.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/ioq3ded"
daemon_user="_ioq3"
daemon_flags="+map q3dm17"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
