# games/ioquake3/pkg/ioq3ded.rc
#!/bin/sh
#
# $OpenBSD: ioq3ded.rc,v 1.2 2016/06/02 11:51:15 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/ioq3ded"
daemon_user="_ioq3"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_bg=YES
rc_reload=NO

rc_cmd $1
