# audio/squeezelite/pkg/squeezelite.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/squeezelite -z"
daemon_user="_squeezelite"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
