# audio/squeezelite/pkg/squeezelite.rc
#!/bin/sh
#
# $OpenBSD: squeezelite.rc,v 1.1 2016/04/04 12:05:37 sthen Exp $

daemon="${TRUEPREFIX}/bin/squeezelite -z"
daemon_user="_squeezelite"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
