# x11/x11vnc/pkg/x11vnc.rc
#!/bin/sh
#
# $OpenBSD: x11vnc.rc,v 1.8 2011/03/17 16:48:47 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/x11vnc"
daemon_flags="-localhost -logfile /var/log/x11vnc -ncache 0 -display :0 -forever -loop100 -auth /etc/X11/xdm/authdir/authfiles/A:0-*"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
