# net/bitlbee/pkg/bitlbee.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/bitlbee -F"
daemon_user="_bitlbee"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
