# sysutils/beats/topbeat/pkg/topbeat.rc
#!/bin/sh
#
# $OpenBSD: topbeat.rc,v 1.1.1.1 2016/04/26 18:33:54 jasper Exp $

daemon="${PREFIX}/bin/topbeat"
daemon_flags="-c ${SYSCONFDIR}/topbeat/topbeat.yml"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
