# mail/milter-greylist/pkg/milter_greylist.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/milter-greylist"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
