# net/openvmps/pkg/vmpsd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/vmpsd"

. /etc/rc.d/rc.subr

rc_cmd $1
