# misc/rocrail/pkg/rocrail.rc
#!/bin/sh
#
# $OpenBSD: rocrail.rc,v 1.3 2017/05/27 13:53:06 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/rocrail"
daemon_user="_rocrail"
daemon_flags="-w /var/rocrail -l ${TRUEPREFIX}/lib/rocrail -img ${TRUEPREFIX}/share/rocrail/images"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
