# misc/rocrail/pkg/rocrail.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/rocrail"
daemon_user="_rocrail"
daemon_flags="-w /var/rocrail -l ${TRUEPREFIX}/lib/rocrail -img ${TRUEPREFIX}/share/rocrail/images"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
