# net/no-ip/pkg/noip2.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/noip2"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
