# mail/mailman/pkg/mailman.rc
#!/bin/sh
#
# $OpenBSD: mailman.rc,v 1.7 2010/12/27 14:50:23 ajacoutot Exp $

daemon="${TRUEPREFIX}/lib/mailman/bin/mailmanctl"
daemon_flags="-s -q"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	su -l -c ${daemon_class} -s ${daemon_shell} ${daemon_user} \
		-c "${daemon} ${daemon_flags} start"
}

rc_stop() {
	${daemon} ${daemon_flags} stop
}

rc_cmd $1
