# mail/mailman/pkg/mailman.rc
#!/bin/sh
#
# $OpenBSD: mailman.rc,v 1.8 2011/01/05 16:05:54 ajacoutot Exp $

daemon="${TRUEPREFIX}/lib/mailman/bin/mailmanctl"
daemon_flags="-s -q"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	${rcexec} "${daemon} ${daemon_flags} start"
}

rc_stop() {
	${daemon} ${daemon_flags} stop
}

rc_cmd $1
