# databases/postgresql/pkg/postgresql.rc
#!/bin/sh
#
# $OpenBSD: postgresql.rc,v 1.7 2011/03/09 18:51:48 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/postgres"
daemon_user="_postgresql"

. /etc/rc.d/rc.subr

bin="${TRUEPREFIX}/bin/pg_ctl"
datadir="/var/postgresql/data"

rc_start() {
	rm -f ${datadir}/postmaster.pid
	${rcexec} "${bin} -D ${datadir} start -l /var/postgresql/logfile"
}

rc_stop() {
	${rcexec} "${bin} -D ${datadir} stop -m fast" >/dev/null || \
		${rcexec} "${bin} -D ${datadir} stop -m immediate"
}

rc_cmd $1
