# net/net-snmp/pkg/netsnmpd.rc
#!/bin/ksh
#
# $OpenBSD: netsnmpd.rc,v 1.5 2018/01/11 19:27:06 rpe Exp $

daemon="${TRUEPREFIX}/sbin/snmpd"

daemon_flags="-u _netsnmp -r"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_cmd $1
