# net/net-snmp/pkg/netsnmptrapd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/snmptrapd"
daemon_flags="-u _netsnmp"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_cmd $1
