# net/net-snmp/pkg/netsnmptrapd.rc
#!/bin/sh
#
# $OpenBSD: netsnmptrapd.rc,v 1.2 2015/11/05 23:43:46 sthen Exp $

daemon="${TRUEPREFIX}/sbin/snmptrapd"
daemon_flags="-u _netsnmp"

. /etc/rc.d/rc.subr

rc_cmd $1
