# net/net-snmp/pkg/netsnmpd.rc
#!/bin/sh
#
# $OpenBSD: netsnmpd.rc,v 1.4 2016/10/11 20:47:43 sthen Exp $

daemon="${TRUEPREFIX}/sbin/snmpd"

daemon_flags="-u _netsnmp -r"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_cmd $1
