# net/net-snmp/pkg/netsnmptrapd.rc
#!/bin/ksh
#
# $OpenBSD: netsnmptrapd.rc,v 1.4 2018/01/11 19:27:06 rpe Exp $

daemon="${TRUEPREFIX}/sbin/snmptrapd"
daemon_flags="-u _netsnmp"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_cmd $1
