# net/net-snmp/pkg/netsnmpd.rc
#!/bin/sh
#
# $OpenBSD: netsnmpd.rc,v 1.1 2011/07/07 20:02:02 sthen Exp $

daemon="${TRUEPREFIX}/sbin/snmpd"

. /etc/rc.d/rc.subr

rc_cmd $1
