# net/net-snmp/pkg/netsnmpd.rc
#!/bin/sh
#
# $OpenBSD: netsnmpd.rc,v 1.3 2016/01/07 12:35:53 sthen Exp $

daemon="${TRUEPREFIX}/sbin/snmpd"

# disable ipv6 mib which needs root; to override the default you may set
# this in rc.conf.local: netsnmpd_flags=" "
daemon_flags="-u _netsnmp -I -ipv6"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_cmd $1
