# net/net-snmp/pkg/netsnmptrapd.rc
#!/bin/sh
#
# $OpenBSD: netsnmptrapd.rc,v 1.3 2016/01/07 12:35:53 sthen Exp $

daemon="${TRUEPREFIX}/sbin/snmptrapd"
daemon_flags="-u _netsnmp"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"

rc_cmd $1
