# net/avahi/pkg/avahi_daemon.rc
#!/bin/sh
#
# $OpenBSD: avahi_daemon.rc,v 1.6 2011/03/17 16:49:27 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/avahi-daemon"
daemon_flags="-D"

. /etc/rc.d/rc.subr

rc_check() {
	${daemon} --check
}

rc_reload() {
	${daemon} --reload
}

rc_stop() {
	${daemon} --kill
}

rc_cmd $1
