# net/avahi/pkg/avahi_daemon.rc
#!/bin/ksh
#
# $OpenBSD: avahi_daemon.rc,v 1.8 2018/01/11 19:27:04 rpe Exp $

daemon="${TRUEPREFIX}/sbin/avahi-daemon -D"

. /etc/rc.d/rc.subr

rc_check() {
	${daemon} --check
}

rc_reload() {
	${daemon} --reload
}

rc_stop() {
	${daemon} --kill
}

rc_cmd $1
