# net/avahi/pkg/avahi_dnsconfd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/avahi-dnsconfd -D"

. /etc/rc.d/rc.subr

rc_check() {
	${daemon} --check
}

rc_reload() {
	${daemon} --refresh
}

rc_stop() {
	${daemon} --kill
}

rc_cmd $1
