# net/avahi/pkg/avahi_daemon.rc
#!/bin/sh
#
# $OpenBSD: avahi_daemon.rc,v 1.2 2010/10/29 12:54:43 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/avahi-daemon"
daemon_flags="-D"

rc_check() {
	${daemon} --check
}

rc_reload() {
	${daemon} --reload
}

rc_stop() {
	${daemon} --kill
}

rc_cmd $1
