# net/avahi/pkg/avahi_daemon.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/avahi-daemon -D"

. /etc/rc.d/rc.subr

rc_check() {
	${daemon} --check
}

rc_reload() {
	${daemon} --reload
}

rc_stop() {
	${daemon} --kill
}

rc_cmd $1
