# net/avahi/pkg/avahi_dnsconfd.rc
#!/bin/sh
#
# $OpenBSD: avahi_dnsconfd.rc,v 1.6 2011/03/17 16:49:27 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/avahi-dnsconfd"
daemon_flags="-D"

. /etc/rc.d/rc.subr

rc_check() {
	${daemon} --check
}

rc_reload() {
	${daemon} --refresh
}

rc_stop() {
	${daemon} --kill
}

rc_cmd $1
