# sysutils/amazon-ssm-agent/pkg/amazon_ssm_agent.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/amazon-ssm-agent"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
