# sysutils/amazon-ssm-agent/pkg/amazon_ssm_agent.rc
#!/bin/sh
#
# $OpenBSD: amazon_ssm_agent.rc,v 1.3 2017/05/27 19:57:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/amazon-ssm-agent"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
