# sysutils/amazon-ssm-agent/pkg/amazon_ssm_agent.rc
#!/bin/sh
#
# $OpenBSD: amazon_ssm_agent.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/amazon-ssm-agent"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
