# comms/lcdproc/pkg/lcdd.rc
#!/bin/sh
#
# $OpenBSD: lcdd.rc,v 1.1 2011/12/02 12:39:23 edd Exp $

daemon="${TRUEPREFIX}/sbin/LCDd"

. /etc/rc.d/rc.subr

rc_cmd $1
