# comms/lcdproc/pkg/lcdexec.rc
#!/bin/ksh
#
# $OpenBSD: lcdexec.rc,v 1.2 2018/01/11 19:27:01 rpe Exp $

daemon="${TRUEPREFIX}/bin/lcdexec"
daemon_user="_lcdproc"

. /etc/rc.d/rc.subr

rc_cmd $1
