# comms/lcdproc/pkg/lcdproc.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/lcdproc"
daemon_user="_lcdproc"

. /etc/rc.d/rc.subr

rc_cmd $1
