# lang/php/5.3/pkg/php_fpm.rc
#!/bin/sh
#
# $OpenBSD: php_fpm.rc,v 1.1 2011/08/27 13:13:13 robert Exp $

daemon="${TRUEPREFIX}/sbin/php-fpm-${PV}"

. /etc/rc.d/rc.subr

pexp="php-fpm-${PV}: master.*"

rc_reload() {
	pkill -USR2 -f "^${pexp}"
}

rc_cmd $1
