# lang/php/5.3/pkg/php_fpm.rc
#!/bin/sh
#
# $OpenBSD: php_fpm.rc,v 1.2 2012/05/19 02:49:45 william Exp $

daemon="${TRUEPREFIX}/sbin/php-fpm-${PV}"

. /etc/rc.d/rc.subr

pexp="php-fpm-${PV}: master.*"
rc_reload=NO

rc_cmd $1
