# www/dansguardian/pkg/dansguardian.rc
#!/bin/sh
#
# $OpenBSD: dansguardian.rc,v 1.1 2013/01/18 15:04:09 sthen Exp $

daemon="${TRUEPREFIX}/sbin/dansguardian"

. /etc/rc.d/rc.subr

rc_cmd $1
