# net/mlvpn/pkg/mlvpn.rc
#!/bin/sh
#
# $OpenBSD: mlvpn.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/mlvpn"
daemon_flags="-c ${SYSCONFDIR}/mlvpn/mlvpn.conf --user _mlvpn"

. /etc/rc.d/rc.subr

pexp="mlvpn: .*priv.*"

rc_cmd $1
