# net/mlvpn/pkg/mlvpn.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/mlvpn"
daemon_flags="-c ${SYSCONFDIR}/mlvpn/mlvpn.conf --user _mlvpn"

. /etc/rc.d/rc.subr

pexp="mlvpn: .*priv.*"
rc_bg=YES

rc_cmd $1
