# sysutils/cfengine/pkg/cf_monitord.rc
#!/bin/ksh

daemon="/var/cfengine/bin/cf-monitord"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
