# sysutils/cfengine/pkg/cf_serverd.rc
#!/bin/ksh
#
# $OpenBSD: cf_serverd.rc,v 1.2 2018/01/11 19:27:10 rpe Exp $

daemon="/var/cfengine/bin/cf-serverd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
