# sysutils/cfengine/pkg/cf_serverd.rc
#!/bin/ksh

daemon="/var/cfengine/bin/cf-serverd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
