# devel/git/pkg/gitdaemon.rc
#!/bin/sh
#
# $OpenBSD: gitdaemon.rc,v 1.1 2016/06/02 18:33:27 jung Exp $

daemon="${TRUEPREFIX}/bin/git daemon --detach"
daemon_flags="--user=_gitdaemon"

. /etc/rc.d/rc.subr

pexp="git-daemon --detach${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
