# devel/git/pkg/gitdaemon.rc
#!/bin/sh
#
# $OpenBSD: gitdaemon.rc,v 1.2 2017/12/01 14:47:45 sthen Exp $

daemon="${TRUEPREFIX}/bin/git daemon --detach"
daemon_flags="--user=_gitdaemon"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/libexec/git/git-daemon --detach${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
