# devel/git/pkg/gitdaemon.rc
#!/bin/ksh
#
# $OpenBSD: gitdaemon.rc,v 1.3 2018/01/11 19:27:02 rpe Exp $

daemon="${TRUEPREFIX}/bin/git daemon --detach"
daemon_flags="--user=_gitdaemon"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/libexec/git/git-daemon --detach${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_cmd $1
