# net/bro/pkg/bro.rc
#!/bin/sh
#
# $OpenBSD: bro.rc,v 1.2 2016/10/07 10:05:16 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/broctl"
daemon_flags="start"

. /etc/rc.d/rc.subr

rc_check() {
	${daemon} status
}

rc_reload() {
	${daemon} update
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
