# net/bro/pkg/zeek.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/zeekctl"
daemon_flags="start"

. /etc/rc.d/rc.subr

rc_check() {
	${daemon} status
}

rc_reload() {
	${daemon} update
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
