# sysutils/ssl_exporter/pkg/ssl_exporter.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/ssl_exporter"
daemon_user="_ssl_exporter"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"
rc_bg=YES
rc_reload=NO

rc_cmd $1
