# audio/mpd/pkg/mpd.rc
#!/bin/sh
#
# $OpenBSD: mpd.rc,v 1.2 2010/10/29 12:54:43 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/mpd"

rc_cmd $1
