# audio/mpd/pkg/mpd.rc
#!/bin/sh

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/mpd"

rc_cmd $1
