# net/ddclient/pkg/ddclient.rc
#!/bin/sh
#
# $OpenBSD: ddclient.rc,v 1.3 2010/12/27 14:50:22 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/ddclient"

. /etc/rc.d/rc.subr

pexp="perl: ddclient"

rc_cmd $1
