# net/ddclient/pkg/ddclient.rc
#!/bin/sh
#
# $OpenBSD: ddclient.rc,v 1.2 2010/12/24 10:40:05 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/sbin/ddclient"

rc_conf

pexp="perl: ddclient"

rc_cmd $1
