# net/ejabberd/pkg/ejabberd.rc
#!/bin/ksh
#
# $OpenBSD: ejabberd.rc,v 1.3 2018/01/11 19:27:05 rpe Exp $

daemon="${TRUEPREFIX}/sbin/ejabberdctl"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_usercheck=NO

rc_start() {
	${rcexec} "${daemon} ${daemon_flags} start"
}

rc_stop() {
	${daemon} ${daemon_flags} stop
}

rc_check() {
	${daemon} ${daemon_flags} status
}

rc_cmd $1
