# sysutils/symon/pkg/symon.rc
#!/bin/ksh
#
# $OpenBSD: symon.rc,v 1.3 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/libexec/symon"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
