# sysutils/symon/pkg/symon.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/symon"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
