# sysutils/symon/pkg/symux.rc
#!/bin/ksh
#
# $OpenBSD: symux.rc,v 1.3 2018/01/11 19:27:10 rpe Exp $

daemon="${TRUEPREFIX}/libexec/symux"

. /etc/rc.d/rc.subr

rc_cmd $1
