# sysutils/symon/pkg/symux.rc
#!/bin/sh
#
# $OpenBSD: symux.rc,v 1.2 2011/05/19 20:01:27 gsoares Exp $

daemon="${TRUEPREFIX}/libexec/symux"

. /etc/rc.d/rc.subr

rc_cmd $1
