# sysutils/symon/pkg/symux.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/libexec/symux"

. /etc/rc.d/rc.subr

rc_cmd $1
