# net/samba/pkg/nmbd.rc
#!/bin/sh
#
# $OpenBSD: nmbd.rc,v 1.6 2012/08/04 15:28:14 ajacoutot Exp $

smbcontrol="${TRUEPREFIX}/bin/smbcontrol"

daemon="${TRUEPREFIX}/libexec/nmbd -D"

. /etc/rc.d/rc.subr

rc_usercheck=NO

rc_check() {
	${smbcontrol} nmbd ping
}

rc_reload() {
	${smbcontrol} nmbd reload-config
}

rc_stop() {
	${smbcontrol} nmbd shutdown
}

rc_cmd $1
