# net/samba/pkg/nmbd.rc
#!/bin/sh
#
# $OpenBSD: nmbd.rc,v 1.4 2011/05/19 20:44:51 gsoares Exp $

daemon="${TRUEPREFIX}/libexec/nmbd"
daemon_flags="-D"

. /etc/rc.d/rc.subr

smbcontrol="${TRUEPREFIX}/bin/smbcontrol"

rc_check() {
	${smbcontrol} nmbd ping
}

rc_reload() {
	${smbcontrol} nmbd reload-config
}

rc_stop() {
	${smbcontrol} nmbd shutdown
}

rc_cmd $1
