# net/samba/pkg/nmbd.rc
#!/bin/sh
#
# $OpenBSD: nmbd.rc,v 1.3 2011/03/17 16:49:27 ajacoutot Exp $

daemon="/usr/local/libexec/nmbd"
daemon_flags="-D"

. /etc/rc.d/rc.subr

smbcontrol="${TRUEPREFIX}/bin/smbcontrol"

rc_check() {
	${smbcontrol} nmbd ping
}

rc_reload() {
	${smbcontrol} nmbd reload-config
}

rc_stop() {
	${smbcontrol} nmbd shutdown
}

rc_cmd $1
