# net/samba/pkg/smbd.rc
#!/bin/sh
#
# $OpenBSD: smbd.rc,v 1.1 2011/01/17 17:58:04 ajacoutot Exp $

daemon="/usr/local/libexec/smbd"
daemon_flags="-D"

. /etc/rc.d/rc.subr

smbcontrol="${TRUEPREFIX}/bin/smbcontrol"

rc_check() {
	${smbcontrol} smbd ping >/dev/null 2>&1
}

rc_reload() {
	${smbcontrol} smbd reload-config
}

rc_stop() {
	${smbcontrol} smbd shutdown 2>/dev/null
}

rc_cmd $1
