# net/samba/pkg/smbd.rc
#!/bin/sh
#
# $OpenBSD: smbd.rc,v 1.4 2011/05/19 20:44:51 gsoares Exp $

daemon="${TRUEPREFIX}/libexec/smbd"
daemon_flags="-D"

. /etc/rc.d/rc.subr

smbcontrol="${TRUEPREFIX}/bin/smbcontrol"

rc_check() {
	${smbcontrol} smbd ping
}

rc_reload() {
	${smbcontrol} smbd reload-config
}

rc_stop() {
	${smbcontrol} smbd shutdown
}

rc_cmd $1
