# net/samba/pkg/smbd.rc
#!/bin/sh
#
# $OpenBSD: smbd.rc,v 1.2 2011/03/09 19:37:47 ajacoutot Exp $

daemon="/usr/local/libexec/smbd"
daemon_flags="-D"

. /etc/rc.d/rc.subr

smbcontrol="${TRUEPREFIX}/bin/smbcontrol"

rc_check() {
	${smbcontrol} smbd ping 2>/dev/null
}

rc_reload() {
	${smbcontrol} smbd reload-config
}

rc_stop() {
	${smbcontrol} smbd shutdown 2>/dev/null
}

rc_cmd $1
