# net/samba/pkg/winbindd.rc
#!/bin/sh
#
# $OpenBSD: winbindd.rc,v 1.4 2011/05/19 20:44:51 gsoares Exp $

daemon="${TRUEPREFIX}/libexec/winbindd"
daemon_flags="-D"

. /etc/rc.d/rc.subr

smbcontrol="${TRUEPREFIX}/bin/smbcontrol"

rc_check() {
	${smbcontrol} winbindd ping
}

rc_reload() {
	${smbcontrol} winbindd reload-config
}

rc_stop() {
	${smbcontrol} winbindd shutdown
}

rc_cmd $1
