# mail/stalwart/pkg/stalwart_mail.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/stalwart-mail"
daemon_logger="daemon.info"
daemon_flags="--config /var/stalwart/etc/config.toml"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
