# www/rt/pkg/rt.rc
#!/bin/sh
#
# $OpenBSD: rt.rc,v 1.2 2011/06/14 13:19:33 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/rt-server"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl -w ${daemon}"
rc_bg=YES
rc_reload=NO

rc_cmd $1
