# www/rt/pkg/rt.rc
#!/bin/ksh
#
# $OpenBSD: rt.rc,v 1.6 2018/01/11 19:27:12 rpe Exp $

daemon="${TRUEPREFIX}/sbin/rt-server"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl -w ${daemon}.fcgi"
rc_bg=YES
rc_reload=NO

rc_cmd $1
