# sysutils/freedt/pkg/svscan.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/svscan"
daemon_flags="/service/"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_post() {
	sleep 1
	${TRUEPREFIX}/bin/svc -X ${daemon_flags}/*
	${TRUEPREFIX}/bin/svc -X ${daemon_flags}/*/log
}

rc_cmd $1
