# productivity/tryton/5.2/trytond/pkg/trytond_cron.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/trytond-cron"
daemon_flags="--config ${SYSCONFDIR}/trytond/server.conf --logconf ${SYSCONFDIR}/trytond/logging.conf"
daemon_user="_trytond"

. /etc/rc.d/rc.subr

pexp="${MODPY_BIN} ${daemon}${daemon_flags:+ ${daemon_flags}}"
rc_bg=YES
rc_reload=NO

rc_cmd $1
