# sysutils/grafana/pkg/grafana.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/grafana-server"
daemon_user="_grafana"
daemon_flags="-homepath ${TRUEPREFIX}/share/grafana -config ${SYSCONFDIR}/grafana/config.ini"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
