# mail/spampd/pkg/spampd.rc
#!/bin/ksh
#
# $OpenBSD: spampd.rc,v 1.4 2021/08/07 10:15:24 sthen Exp $

daemon="${TRUEPREFIX}/sbin/spampd"
daemon_flags="--a --p=/var/spampd/spampd.pid"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl -T ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
