# mail/spampd/pkg/spampd.rc
#!/bin/ksh
#
# $OpenBSD: spampd.rc,v 1.3 2018/01/11 19:27:04 rpe Exp $

daemon="${TRUEPREFIX}/sbin/spampd"
daemon_flags="--a --aw --p=/var/spampd/spampd.pid"

. /etc/rc.d/rc.subr

pexp="/usr/bin/perl -T ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
