# www/havp/pkg/havp.rc
#!/bin/sh
#
# $OpenBSD: havp.rc,v 1.2 2011/11/30 11:00:59 robert Exp $

daemon="${TRUEPREFIX}/sbin/havp"

. /etc/rc.d/rc.subr

rc_pre() {
	install -d -o _havp -g _havp /var/run/havp/ /var/tmp/havp/
}

rc_cmd $1
