# security/snort2pf/pkg/snort2pf.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/snort2pf"

. /etc/rc.d/rc.subr

pexp="perl: snort2pf .*"

rc_bg=YES
rc_reload=NO

rc_cmd $1
