# security/snort2pf/pkg/snort2pf.rc
#!/bin/sh
#
# $OpenBSD: snort2pf.rc,v 1.3 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/sbin/snort2pf"

. /etc/rc.d/rc.subr

pexp="perl: snort2pf .*"

rc_reload=NO

rc_cmd $1
