# net/pure-ftpd/pkg/pure_ftpd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/pure-ftpd"
daemon_flags="-A -B -H -u1000"

. /etc/rc.d/rc.subr

pexp="pure-ftpd: -pure-ftpd \(SERVER\)"
rc_reload=NO

rc_cmd $1
