# security/dropbear/pkg/dropbear.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/dropbear"
daemon_flags="-p 127.0.0.1:8022 -R"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
