# security/dropbear/pkg/dropbear.rc
#!/bin/ksh
#
# $OpenBSD: dropbear.rc,v 1.1 2018/03/13 03:50:48 bket Exp $

daemon="${TRUEPREFIX}/sbin/dropbear"
daemon_flags="-p 127.0.0.1:8022 -R"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
