# www/polipo/pkg/polipo.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/polipo daemonise=yes"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
