# www/polipo/pkg/polipo.rc
#!/bin/ksh
#
# $OpenBSD: polipo.rc,v 1.7 2018/01/11 19:27:12 rpe Exp $

daemon="${TRUEPREFIX}/bin/polipo daemonise=yes"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
