# www/polipo/pkg/polipo.rc
#!/bin/sh
#
# $OpenBSD: polipo.rc,v 1.4 2010/12/24 10:40:06 ajacoutot Exp $

. /etc/rc.d/rc.subr

daemon="${TRUEPREFIX}/bin/polipo"
daemon_flags="daemonise=yes"

rc_conf

rc_reload=NO

rc_cmd $1
