# net/tinc/pkg/tincd.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/tincd"
daemon_flags="-U _tinc --chroot --net=example"

. /etc/rc.d/rc.subr

rc_cmd $1
