# sysutils/rundeck/pkg/rundeck.rc
#!/bin/ksh

daemon=$(${LOCALBASE}/bin/javaPathHelper -c rundeck)
daemon_flags="-XX:MaxPermSize=256m -Xmx1024m -Djava.io.tmpdir=/var/cache/rundeck -jar ${RDECK_BASE}/rundeck-launcher.jar -b ${RDECK_BASE}"
daemon_user="_rundeck"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO
rc_stop_signal=KILL

rc_cmd $1
