# security/opendnssec/pkg/opendnssec.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/sbin/ods-control"

. /etc/rc.d/rc.subr

rc_reload=NO

pexp="${TRUEPREFIX}/sbin/ods-(enforcerd|signerd)"

rc_pre() {
	install -d -o _opendnssec /var/run/opendnssec/
}

rc_start() {
	${rcexec} "${daemon} start"
}

rc_stop() {
	${daemon} stop
}

rc_cmd $1
