# net/go-ipfs/pkg/go_ipfs.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/ipfs daemon"
daemon_user="_go-ipfs"
daemon_timeout=60

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_start() {
	${rcexec} "IPFS_PATH=\"${VARBASE}/go-ipfs\" ${daemon} ${daemon_flags}"
}

rc_cmd $1
