# net/ices2/pkg/ices2.rc
#!/bin/sh
#
# $OpenBSD: ices2.rc,v 1.2 2017/05/27 13:53:07 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/ices2"
daemon_flags="${SYSCONFDIR}/ices2.xml"

. /etc/rc.d/rc.subr

rc_cmd $1
