# net/apache-activemq/pkg/activemq.rc
#!/bin/sh
#
# $OpenBSD: activemq.rc,v 1.2 2013/03/04 15:48:41 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/activemq"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_start() {
	${rcexec} "PATH=$PATH:$(dirname $(${LOCALBASE}/bin/javaPathHelper -c apache-activemq)) ${daemon} start"
}
 
rc_check() {
	PATH=$PATH:$(dirname $(${LOCALBASE}/bin/javaPathHelper -c apache-activemq)) ${daemon} status   
}

rc_stop() {
	pid=$(rc_check | tail -1 | perl -pi -e 's,.*pid \x27(.*)\x27.*,$1,g')
	kill ${pid}
}

rc_cmd $1
