# mail/archiveopteryx/pkg/aox.rc
#!/bin/sh
#
# $OpenBSD: aox.rc,v 1.2 2011/03/09 18:51:48 ajacoutot Exp $

daemon="${TRUEPREFIX}/bin/aox"

. /etc/rc.d/rc.subr

pexp="${TRUEPREFIX}/sbin/archiveopteryx -f"
rc_reload=NO

rc_start() {
	${daemon} start 2>/dev/null
}

rc_stop() {
	${daemon} stop 2>/dev/null
}

rc_cmd $1
