# net/dnsfilter/pkg/dnsfilter.rc
#!/bin/sh
#
# $OpenBSD: dnsfilter.rc,v 1.1.1.1 2012/12/04 16:10:25 florian Exp $

daemon="${TRUEPREFIX}/sbin/dnsfilter"

. /etc/rc.d/rc.subr

pexp="dnsfilter: parent"
rc_reload=NO

rc_cmd $1
