# audio/shairport-sync/pkg/shairport_sync.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/shairport-sync"
daemon_user="_shairport"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_cmd $1
